# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/file_event_win_ripzip_attack.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Potential RipZip Attack on Startup Folder</title><id type="str">a6976974-ea6f-4e97-818e-ea08625c52cb</id><status type="str">experimental</status><description type="str">Detects a phishing attack which expands a ZIP file containing a malicious shortcut.
If the victim expands the ZIP file via the explorer process, then the explorer process expands the malicious ZIP file and drops a malicious shortcut redirected to a backdoor into the Startup folder.
Additionally, the file name of the malicious shortcut in Startup folder contains {0AFACED1-E828-11D1-9187-B532F1E9575D} meaning the folder shortcut operation.
</description><references type="list"><item type="str">https://twitter.com/jonasLyk/status/1549338335243534336?t=CrmPocBGLbDyE4p6zTX1cg&amp;s=19</item></references><author type="str">Greg (rule)</author><date type="str">2022/07/21</date><modified type="str">2023/01/05</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1547</item></tags><logsource type="dict"><category type="str">file_event</category><product type="str">windows</product></logsource><detection type="dict"><selection type="dict"><key name="TargetFilename|contains|all" type="list"><item type="str">\Microsoft\Windows\Start Menu\Programs\Startup</item><item type="str">.lnk.{0AFACED1-E828-11D1-9187-B532F1E9575D}</item></key><key name="Image|endswith" type="str">\explorer.exe</key></selection><condition type="str">selection</condition></detection><falsepositives type="list"><item type="str">Unknown</item></falsepositives><level type="str">high</level></root>