# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/file_event_lnx_persistence_cron_files.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Persistence Via Cron Files</title><id type="str">6c4e2f43-d94d-4ead-b64d-97e53fa2bd05</id><status type="str">test</status><description type="str">Detects creation of cron file or files in Cron directories which could indicates potential persistence.</description><references type="list"><item type="str">https://github.com/microsoft/MSTIC-Sysmon/blob/f1477c0512b0747c1455283069c21faec758e29d/linux/configs/attack-based/persistence/T1053.003_Cron_Activity.xml</item></references><author type="str">Roberto Rodriguez (Cyb3rWard0g), OTR (Open Threat Research), MSTIC</author><date type="str">2021/10/15</date><modified type="str">2022/12/31</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1053.003</item></tags><logsource type="dict"><product type="str">linux</product><category type="str">file_event</category></logsource><detection type="dict"><selection1 type="dict"><key name="TargetFilename|startswith" type="list"><item type="str">/etc/cron.d/</item><item type="str">/etc/cron.daily/</item><item type="str">/etc/cron.hourly/</item><item type="str">/etc/cron.monthly/</item><item type="str">/etc/cron.weekly/</item><item type="str">/var/spool/cron/crontabs/</item></key></selection1><selection2 type="dict"><key name="TargetFilename|contains" type="list"><item type="str">/etc/cron.allow</item><item type="str">/etc/cron.deny</item><item type="str">/etc/crontab</item></key></selection2><condition type="str">1 of selection*</condition></detection><falsepositives type="list"><item type="str">Any legitimate cron file.</item></falsepositives><level type="str">medium</level></root>