# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/file_event_win_susp_startup_folder_persistence.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Suspicious Startup Folder Persistence</title><id type="str">28208707-fe31-437f-9a7f-4b1108b94d2e</id><related type="list"><item type="dict"><id type="str">2aa0a6b4-a865-495b-ab51-c28249537b75</id><type type="str">similar</type></item></related><status type="str">experimental</status><description type="str">Detects when a file with a suspicious extension is created in the startup folder</description><references type="list"><item type="str">https://github.com/last-byte/PersistenceSniper</item></references><author type="str">Nasreddine Bencherchali (Nextron Systems)</author><date type="str">2022/08/10</date><modified type="str">2023/01/06</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1547.001</item></tags><logsource type="dict"><product type="str">windows</product><category type="str">file_event</category></logsource><detection type="dict"><selection type="dict"><key name="TargetFilename|contains" type="str">\Windows\Start Menu\Programs\Startup\</key><key name="TargetFilename|endswith" type="list"><item type="str">.vbs</item><item type="str">.vbe</item><item type="str">.bat</item><item type="str">.ps1</item><item type="str">.hta</item><item type="str">.dll</item><item type="str">.jar</item><item type="str">.msi</item><item type="str">.scr</item><item type="str">.cmd</item></key></selection><condition type="str">selection</condition></detection><falsepositives type="list"><item type="str">Rare legitimate usage of some of the extensions mentioned in the rule</item></falsepositives><level type="str">high</level></root>