# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/file_event_win_msdt_susp_directories.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">File Creation In Suspicious Directory By Msdt.EXE</title><id type="str">318557a5-150c-4c8d-b70e-a9910e199857</id><status type="str">experimental</status><description type="str">Detects msdt.exe creating files in suspicious directories which could be a sign of exploitation of either Follina or Dogwalk vulnerabilities</description><references type="list"><item type="str">https://irsl.medium.com/the-trouble-with-microsofts-troubleshooters-6e32fc80b8bd</item><item type="str">https://msrc-blog.microsoft.com/2022/05/30/guidance-for-cve-2022-30190-microsoft-support-diagnostic-tool-vulnerability/</item></references><author type="str">Vadim Varganov, Florian Roth (Nextron Systems)</author><date type="str">2022/08/24</date><modified type="str">2023/02/23</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1547.001</item><item type="str">cve.2022.30190</item></tags><logsource type="dict"><category type="str">file_event</category><product type="str">windows</product></logsource><detection type="dict"><selection type="dict"><key name="Image|endswith" type="str">\msdt.exe</key><key name="TargetFilename|contains" type="list"><item type="str">\Desktop\</item><item type="str">\Start Menu\Programs\Startup\</item><item type="str">C:\PerfLogs\</item><item type="str">C:\ProgramData\</item><item type="str">C:\Users\Public\</item></key></selection><condition type="str">selection</condition></detection><falsepositives type="list"><item type="str">Unknown</item></falsepositives><level type="str">high</level></root>