# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/lnx_auditd_systemd_service_creation.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Systemd Service Creation</title><id type="str">1bac86ba-41aa-4f62-9d6b-405eac99b485</id><status type="str">test</status><description type="str">Detects a creation of systemd services which could be used by adversaries to execute malicious code.</description><references type="list"><item type="str">https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1543.002/T1543.002.md</item></references><author type="str">Pawel Mazur</author><date type="str">2022/02/03</date><modified type="str">2022/02/06</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1543.002</item></tags><logsource type="dict"><product type="str">linux</product><service type="str">auditd</service></logsource><detection type="dict"><path type="dict"><type type="str">PATH</type><nametype type="str">CREATE</nametype></path><name_1 type="dict"><key name="name|startswith" type="list"><item type="str">/usr/lib/systemd/system/</item><item type="str">/etc/systemd/system/</item></key></name_1><name_2 type="dict"><key name="name|contains" type="str">/.config/systemd/user/</key></name_2><condition type="str">path and 1 of name_*</condition></detection><falsepositives type="list"><item type="str">Admin work like legit service installs.</item></falsepositives><level type="str">medium</level></root>