# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/posh_ps_cmdlet_scheduled_task.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Powershell Create Scheduled Task</title><id type="str">363eccc0-279a-4ccf-a3ab-24c2e63b11fb</id><status type="str">test</status><description type="str">Adversaries may abuse the Windows Task Scheduler to perform task scheduling for initial or recurring execution of malicious code</description><references type="list"><item type="str">https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1053.005/T1053.005.md#atomic-test-4---powershell-cmdlet-scheduled-task</item><item type="str">https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1053.005/T1053.005.md#atomic-test-6---wmi-invoke-cimmethod-scheduled-task</item></references><author type="str">frack113</author><date type="str">2021/12/28</date><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1053.005</item></tags><logsource type="dict"><product type="str">windows</product><category type="str">ps_script</category><definition type="str">Requirements: Script Block Logging must be enabled</definition></logsource><detection type="dict"><selection_cmdlet type="dict"><key name="ScriptBlockText|contains" type="list"><item type="str">New-ScheduledTaskAction</item><item type="str">New-ScheduledTaskTrigger</item><item type="str">New-ScheduledTaskPrincipal</item><item type="str">New-ScheduledTaskSettingsSet</item><item type="str">New-ScheduledTask</item><item type="str">Register-ScheduledTask</item></key></selection_cmdlet><selection_cimmethod type="dict"><key name="ScriptBlockText|contains|all" type="list"><item type="str">Invoke-CimMethod</item><item type="str">-ClassName</item><item type="str">PS_ScheduledTask</item><item type="str">-NameSpace</item><item type="str">Root\Microsoft\Windows\TaskScheduler</item></key></selection_cimmethod><condition type="str">1 of selection_*</condition></detection><falsepositives type="list"><item type="str">Unknown</item></falsepositives><level type="str">medium</level></root>