# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_win_apt_actinium_persistence.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Potential ACTINIUM Persistence Activity</title><id type="str">e1118a8f-82f5-44b3-bb6b-8a284e5df602</id><status type="str">test</status><description type="str">Detects specific process parameters as used by ACTINIUM scheduled task persistence creation.</description><references type="list"><item type="str">https://www.microsoft.com/security/blog/2022/02/04/actinium-targets-ukrainian-organizations</item></references><author type="str">Andreas Hunkeler (@Karneades)</author><date type="str">2022/02/07</date><modified type="str">2023/03/18</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1053</item><item type="str">attack.t1053.005</item><item type="str">detection.emerging_threats</item></tags><logsource type="dict"><category type="str">process_creation</category><product type="str">windows</product></logsource><detection type="dict"><selection type="dict"><key name="CommandLine|contains|all" type="list"><item type="str">schtasks</item><item type="str">create</item><item type="str">wscript</item><item type="str"> /e:vbscript</item></key></selection><condition type="str">selection</condition></detection><falsepositives type="list"><item type="str">Unlikely</item></falsepositives><level type="str">high</level></root>