# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_win_apt_oilrig_mar18.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">OilRig APT Activity</title><id type="str">ce6e34ca-966d-41c9-8d93-5b06c8b97a06</id><related type="list"><item type="dict"><id type="str">53ba33fd-3a50-4468-a5ef-c583635cfa92</id><type type="str">similar</type></item><item type="dict"><id type="str">c0580559-a6bd-4ef6-b9b7-83703d98b561</id><type type="str">similar</type></item><item type="dict"><id type="str">7bdf2a7c-3acc-4091-9581-0a77dad1c5b5</id><type type="str">similar</type></item></related><status type="str">test</status><description type="str">Detects OilRig activity as reported by Nyotron in their March 2018 report</description><references type="list"><item type="str">https://web.archive.org/web/20180402134442/https://nyotron.com/wp-content/uploads/2018/03/Nyotron-OilRig-Malware-Report-March-2018C.pdf</item></references><author type="str">Florian Roth (Nextron Systems), Markus Neis, Jonhnathan Ribeiro, Daniil Yugoslavskiy, oscd.community</author><date type="str">2018/03/23</date><modified type="str">2023/03/08</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.g0049</item><item type="str">attack.t1053.005</item><item type="str">attack.s0111</item><item type="str">attack.t1543.003</item><item type="str">attack.defense_evasion</item><item type="str">attack.t1112</item><item type="str">attack.command_and_control</item><item type="str">attack.t1071.004</item><item type="str">detection.emerging_threats</item></tags><logsource type="dict"><category type="str">process_creation</category><product type="str">windows</product></logsource><detection type="dict"><selection_schtasks type="dict"><key name="CommandLine|contains|all" type="list"><item type="str">SC Scheduled Scan</item><item type="str">\microsoft\Taskbar\autoit3.exe</item></key></selection_schtasks><selection_temp type="dict"><key name="Image|contains" type="str">\Windows\Temp\DB\</key><key name="Image|endswith" type="str">.exe</key></selection_temp><selection_service type="dict"><Image type="str">C:\Windows\system32\Service.exe</Image><key name="CommandLine|contains" type="list"><item type="str">i</item><item type="str">u</item></key></selection_service><selection_autoit type="dict"><key name="ParentImage|endswith" type="str">\local\microsoft\Taskbar\autoit3.exe</key><key name="CommandLine|contains|all" type="list"><item type="str">nslookup.exe</item><item type="str">-q=TXT</item></key></selection_autoit><condition type="str">1 of selection_*</condition></detection><falsepositives type="list"><item type="str">Unlikely</item></falsepositives><level type="str">critical</level></root>