# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_win_schtasks_creation.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Scheduled Task Creation</title><id type="str">92626ddd-662c-49e3-ac59-f6535f12d189</id><status type="str">test</status><description type="str">Detects the creation of scheduled tasks in user session</description><author type="str">Florian Roth (Nextron Systems)</author><date type="str">2019/01/16</date><modified type="str">2022/10/09</modified><tags type="list"><item type="str">attack.execution</item><item type="str">attack.persistence</item><item type="str">attack.privilege_escalation</item><item type="str">attack.t1053.005</item><item type="str">attack.s0111</item><item type="str">car.2013-08-001</item></tags><logsource type="dict"><category type="str">process_creation</category><product type="str">windows</product></logsource><detection type="dict"><selection type="dict"><key name="Image|endswith" type="str">\schtasks.exe</key><key name="CommandLine|contains" type="str"> /create </key></selection><filter type="dict"><key name="User|contains" type="list"><item type="str">AUTHORI</item><item type="str">AUTORI</item></key></filter><condition type="str">selection and not filter</condition></detection><fields type="list"><item type="str">CommandLine</item><item type="str">ParentCommandLine</item></fields><falsepositives type="list"><item type="str">Administrative activity</item><item type="str">Software installation</item></falsepositives><level type="str">low</level></root>