# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_win_schtasks_guid_task_name.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Suspicious Scheduled Task Name As GUID</title><id type="str">ff2fff64-4cd6-4a2b-ba7d-e28a30bbe66b</id><status type="str">experimental</status><description type="str">Detects creation of a scheduled task with a GUID like name</description><references type="list"><item type="str">https://thedfirreport.com/2022/10/31/follina-exploit-leads-to-domain-compromise/</item><item type="str">https://thedfirreport.com/2022/02/21/qbot-and-zerologon-lead-to-full-domain-compromise/</item></references><author type="str">Nasreddine Bencherchali (Nextron Systems)</author><date type="str">2022/10/31</date><tags type="list"><item type="str">attack.execution</item><item type="str">attack.t1053.005</item></tags><logsource type="dict"><product type="str">windows</product><category type="str">process_creation</category></logsource><detection type="dict"><selection_img type="dict"><key name="Image|endswith" type="str">\schtasks.exe</key><key name="CommandLine|contains" type="str">/Create </key></selection_img><selection_tn type="dict"><key name="CommandLine|contains" type="list"><item type="str">/TN &quot;{</item><item type="str">/TN &apos;{</item><item type="str">/TN {</item></key></selection_tn><selection_end type="dict"><key name="CommandLine|contains" type="list"><item type="str">}&quot;</item><item type="str">}&apos;</item><item type="str">} </item></key></selection_end><condition type="str">all of selection_*</condition></detection><falsepositives type="list"><item type="str">Legitimate software naming their tasks as GUIDs</item></falsepositives><level type="str">medium</level></root>