# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_win_schtasks_appdata_local_system.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Suspicious Schtasks Execution AppData Folder</title><id type="str">c5c00f49-b3f9-45a6-997e-cfdecc6e1967</id><status type="str">experimental</status><description type="str">Detects the creation of a schtask that executes a file from C:\Users\&lt;USER&gt;\AppData\Local</description><references type="list"><item type="str">https://thedfirreport.com/2022/02/21/qbot-and-zerologon-lead-to-full-domain-compromise/</item></references><author type="str">pH-T (Nextron Systems), Nasreddine Bencherchali (Nextron Systems)</author><date type="str">2022/03/15</date><modified type="str">2022/07/28</modified><tags type="list"><item type="str">attack.execution</item><item type="str">attack.persistence</item><item type="str">attack.t1053.005</item><item type="str">attack.t1059.001</item></tags><logsource type="dict"><product type="str">windows</product><category type="str">process_creation</category></logsource><detection type="dict"><selection type="dict"><key name="Image|endswith" type="str">\schtasks.exe</key><key name="CommandLine|contains|all" type="list"><item type="str">/Create</item><item type="str">/RU</item><item type="str">/TR</item><item type="str">C:\Users\</item><item type="str">\AppData\Local\</item></key><key name="CommandLine|contains" type="list"><item type="str">NT AUT</item><item type="str"> SYSTEM </item></key></selection><filter type="dict"><key name="ParentImage|contains|all" type="list"><item type="str">\AppData\Local\Temp\</item><item type="str">TeamViewer_.exe</item></key><key name="Image|endswith" type="str">\schtasks.exe</key><key name="CommandLine|contains" type="str">/TN TVInstallRestore</key></filter><condition type="str">selection and not filter</condition></detection><falsepositives type="list"><item type="str">Unknown</item></falsepositives><level type="str">high</level></root>