# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_win_reg_add_run_key.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Potential Persistence Attempt Via Run Keys Using Reg.EXE</title><id type="str">de587dce-915e-4218-aac4-835ca6af6f70</id><status type="str">test</status><description type="str">Detects suspicious command line reg.exe tool adding key to RUN key in Registry</description><references type="list"><item type="str">https://app.any.run/tasks/9c0f37bc-867a-4314-b685-e101566766d7/</item><item type="str">https://docs.microsoft.com/en-us/windows/win32/setupapi/run-and-runonce-registry-keys</item></references><author type="str">Florian Roth (Nextron Systems)</author><date type="str">2021/06/28</date><modified type="str">2023/01/30</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1547.001</item></tags><logsource type="dict"><category type="str">process_creation</category><product type="str">windows</product></logsource><detection type="dict"><selection type="dict"><key name="CommandLine|contains|all" type="list"><item type="str">reg</item><item type="str"> ADD </item><item type="str">Software\Microsoft\Windows\CurrentVersion\Run</item></key></selection><condition type="str">selection</condition></detection><falsepositives type="list"><item type="str">Legitimate software automatically (mostly, during installation) sets up autorun keys for legitimate reasons.</item><item type="str">Legitimate administrator sets up autorun keys for legitimate reasons.</item><item type="str">Discord</item></falsepositives><level type="str">medium</level></root>