# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_macos_schedule_task_job_cron.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Scheduled Cron Task/Job - MacOs</title><id type="str">7c3b43d8-d794-47d2-800a-d277715aa460</id><status type="str">test</status><description type="str">Detects abuse of the cron utility to perform task scheduling for initial or recurring execution of malicious code. Detection will focus on crontab jobs uploaded from the tmp folder.</description><references type="list"><item type="str">https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1053.003/T1053.003.md</item></references><author type="str">Alejandro Ortuno, oscd.community</author><date type="str">2020/10/06</date><modified type="str">2022/11/27</modified><tags type="list"><item type="str">attack.execution</item><item type="str">attack.persistence</item><item type="str">attack.privilege_escalation</item><item type="str">attack.t1053.003</item></tags><logsource type="dict"><category type="str">process_creation</category><product type="str">macos</product></logsource><detection type="dict"><selection type="dict"><key name="Image|endswith" type="str">/crontab</key><key name="CommandLine|contains" type="str">/tmp/</key></selection><condition type="str">selection</condition></detection><falsepositives type="list"><item type="str">Legitimate administration activities</item></falsepositives><level type="str">medium</level></root>