# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_macos_persistence_via_plistbuddy.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Potential Persistence Via PlistBuddy</title><id type="str">65d506d3-fcfe-4071-b4b2-bcefe721bbbb</id><status type="str">experimental</status><description type="str">Detects potential persistence activity using LaunchAgents or LaunchDaemons via the PlistBuddy utility</description><references type="list"><item type="str">https://redcanary.com/blog/clipping-silver-sparrows-wings/</item><item type="str">https://www.manpagez.com/man/8/PlistBuddy/</item></references><author type="str">Sohan G (D4rkCiph3r)</author><date type="str">2023/02/18</date><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1543.001</item><item type="str">attack.t1543.004</item></tags><logsource type="dict"><category type="str">process_creation</category><product type="str">macos</product></logsource><detection type="dict"><selection type="dict"><key name="Image|endswith" type="str">/PlistBuddy</key><key name="CommandLine|contains|all" type="list"><item type="str">RunAtLoad</item><item type="str">true</item></key><key name="CommandLine|contains" type="list"><item type="str">LaunchAgents</item><item type="str">LaunchDaemons</item></key></selection><condition type="str">selection</condition></detection><falsepositives type="list"><item type="str">Unknown</item></falsepositives><level type="str">high</level></root>