# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_win_apt_slingshot.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Defrag Deactivation</title><id type="str">958d81aa-8566-4cea-a565-59ccd4df27b0</id><status type="str">test</status><description type="str">Detects the deactivation and disabling of the Scheduled defragmentation task as seen by Slingshot APT group</description><references type="list"><item type="str">https://securelist.com/apt-slingshot/84312/</item></references><author type="str">Florian Roth (Nextron Systems), Bartlomiej Czyz (@bczyz1)</author><date type="str">2019/03/04</date><modified type="str">2022/10/09</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1053.005</item><item type="str">attack.s0111</item><item type="str">detection.emerging_threats</item></tags><logsource type="dict"><category type="str">process_creation</category><product type="str">windows</product></logsource><detection type="dict"><selection type="dict"><key name="Image|endswith" type="str">\schtasks.exe</key><key name="CommandLine|contains" type="list"><item type="str">/delete</item><item type="str">/change</item></key><key name="CommandLine|contains|all" type="list"><item type="str">/TN</item><item type="str">\Microsoft\Windows\Defrag\ScheduledDefrag</item></key></selection><condition type="str">selection</condition></detection><falsepositives type="list"><item type="str">Unknown</item></falsepositives><level type="str">medium</level></root>