# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_win_schtasks_parent.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Suspicious Add Scheduled Task Parent</title><id type="str">9494479d-d994-40bf-a8b1-eea890237021</id><status type="str">experimental</status><description type="str">Detects suspicious scheduled task creations from a parent stored in a temporary folder</description><references type="list"><item type="str">https://app.any.run/tasks/649e7b46-9bec-4d05-98a5-dfa9a13eaae5/</item></references><author type="str">Florian Roth (Nextron Systems)</author><date type="str">2022/02/23</date><modified type="str">2022/06/02</modified><tags type="list"><item type="str">attack.execution</item><item type="str">attack.t1053.005</item></tags><logsource type="dict"><product type="str">windows</product><category type="str">process_creation</category></logsource><detection type="dict"><selection type="dict"><key name="Image|endswith" type="str">\schtasks.exe</key><key name="CommandLine|contains" type="str">/Create </key><key name="ParentImage|contains" type="list"><item type="str">\AppData\Local\</item><item type="str">\AppData\Roaming\</item><item type="str">\Temporary Internet</item><item type="str">\Users\Public\</item></key></selection><filter type="dict"><key name="CommandLine|contains" type="list"><item type="str">update_task.xml</item><item type="str">unattended.ini</item></key></filter><condition type="str">selection and not 1 of filter*</condition></detection><falsepositives type="list"><item type="str">Software installers that run from temporary folders and also install scheduled tasks</item></falsepositives><level type="str">medium</level></root>