# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_win_schtasks_creation_temp_folder.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Suspicious Scheduled Task Creation Involving Temp Folder</title><id type="str">39019a4e-317f-4ce3-ae63-309a8c6b53c5</id><status type="str">test</status><description type="str">Detects the creation of scheduled tasks that involves a temporary folder and runs only once</description><references type="list"><item type="str">https://discuss.elastic.co/t/detection-and-response-for-hafnium-activity/266289/3</item></references><author type="str">Florian Roth (Nextron Systems)</author><date type="str">2021/03/11</date><modified type="str">2022/10/09</modified><tags type="list"><item type="str">attack.execution</item><item type="str">attack.persistence</item><item type="str">attack.t1053.005</item></tags><logsource type="dict"><category type="str">process_creation</category><product type="str">windows</product></logsource><detection type="dict"><selection type="dict"><key name="Image|endswith" type="str">\schtasks.exe</key><key name="CommandLine|contains|all" type="list"><item type="str"> /create </item><item type="str"> /sc once </item><item type="str">\Temp\</item></key></selection><condition type="str">selection</condition></detection><fields type="list"><item type="str">CommandLine</item><item type="str">ParentCommandLine</item></fields><falsepositives type="list"><item type="str">Administrative activity</item><item type="str">Software installation</item></falsepositives><level type="str">high</level></root>