# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_win_apt_unc2452_vbscript_pattern.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Suspicious VBScript UN2452 Pattern</title><id type="str">20c3f09d-c53d-4e85-8b74-6aa50e2f1b61</id><status type="str">test</status><description type="str">Detects suspicious inline VBScript keywords as used by UNC2452</description><references type="list"><item type="str">https://www.microsoft.com/security/blog/2021/03/04/goldmax-goldfinder-sibot-analyzing-nobelium-malware/</item></references><author type="str">Florian Roth (Nextron Systems)</author><date type="str">2021/03/05</date><modified type="str">2022/10/09</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1547.001</item><item type="str">detection.emerging_threats</item></tags><logsource type="dict"><category type="str">process_creation</category><product type="str">windows</product></logsource><detection type="dict"><selection type="dict"><key name="CommandLine|contains|all" type="list"><item type="str">Execute</item><item type="str">CreateObject</item><item type="str">RegRead</item><item type="str">window.close</item><item type="str">\Microsoft\Windows\CurrentVersion</item></key></selection><filter type="dict"><key name="CommandLine|contains" type="str">\Software\Microsoft\Windows\CurrentVersion\Run</key></filter><condition type="str">selection and not filter</condition></detection><falsepositives type="list"><item type="str">Unknown</item></falsepositives><level type="str">high</level></root>