# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_win_apt_hafnium.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">HAFNIUM Exchange Exploitation Activity</title><id type="str">bbb2dedd-a0e3-46ab-ba6c-6c82ae7a9aa7</id><status type="str">test</status><description type="str">Detects activity observed by different researchers to be HAFNIUM group activity (or related) on Exchange servers</description><references type="list"><item type="str">https://blog.truesec.com/2021/03/07/exchange-zero-day-proxylogon-and-hafnium/</item><item type="str">https://www.microsoft.com/security/blog/2021/03/02/hafnium-targeting-exchange-servers/</item><item type="str">https://discuss.elastic.co/t/detection-and-response-for-hafnium-activity/266289/3</item><item type="str">https://twitter.com/GadixCRK/status/1369313704869834753?s=20</item><item type="str">https://twitter.com/BleepinComputer/status/1372218235949617161</item></references><author type="str">Florian Roth (Nextron Systems)</author><date type="str">2021/03/09</date><modified type="str">2023/03/09</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1546</item><item type="str">attack.t1053</item><item type="str">detection.emerging_threats</item></tags><logsource type="dict"><category type="str">process_creation</category><product type="str">windows</product></logsource><detection type="dict"><selection_attrib type="dict"><key name="CommandLine|contains|all" type="list"><item type="str">attrib</item><item type="str"> +h </item><item type="str"> +s </item><item type="str"> +r </item><item type="str">.aspx</item></key></selection_attrib><selection_vsperfmon type="list"><item type="dict"><key name="Image|contains" type="str">\ProgramData\VSPerfMon\</key></item><item type="dict"><key name="CommandLine|contains|all" type="list"><item type="str">schtasks</item><item type="str">VSPerfMon</item></key></item></selection_vsperfmon><selection_opera_1 type="dict"><key name="Image|endswith" type="str">Opera_browser.exe</key><key name="ParentImage|endswith" type="list"><item type="str">\services.exe</item><item type="str">\svchost.exe</item></key></selection_opera_1><selection_opera_2 type="dict"><key name="Image|endswith" type="str">Users\Public\opera\Opera_browser.exe</key></selection_opera_2><selection_vssadmin type="dict"><key name="CommandLine|contains|all" type="list"><item type="str">vssadmin list shadows</item><item type="str">Temp\__output</item></key></selection_vssadmin><selection_makecab_1 type="dict"><key name="Image|endswith" type="str">\makecab.exe</key><key name="CommandLine|contains|all" type="list"><item type="str">inetpub\wwwroot\</item><item type="str">.dmp.zip</item></key></selection_makecab_1><selection_makecab_2 type="dict"><key name="Image|endswith" type="str">\makecab.exe</key><key name="CommandLine|contains" type="list"><item type="str">Microsoft\Exchange Server\</item><item type="str">compressionmemory</item><item type="str">.gif</item></key></selection_makecab_2><selection_7zip type="dict"><key name="CommandLine|contains|all" type="list"><item type="str"> -t7z </item><item type="str">C:\Programdata\pst</item><item type="str">\it.zip</item></key></selection_7zip><selection_rundll32 type="dict"><key name="CommandLine|contains|all" type="list"><item type="str">\comsvcs.dll</item><item type="str">Minidump</item><item type="str">full </item><item type="str">\inetpub\wwwroot</item></key></selection_rundll32><selection_other type="dict"><key name="CommandLine|contains" type="list"><item type="str">Windows\Temp\xx.bat</item><item type="str">Windows\WwanSvcdcs</item><item type="str">Windows\Temp\cw.exe</item></key></selection_other><condition type="str">1 of selection*</condition></detection><falsepositives type="list"><item type="str">Unlikely</item></falsepositives><level type="str">critical</level></root>