# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_win_sc_service_path_modification.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Suspicious Service Path Modification</title><id type="str">138d3531-8793-4f50-a2cd-f291b2863d78</id><status type="str">test</status><description type="str">Detects service path modification via the &quot;sc&quot; binary to a suspicious command or path</description><references type="list"><item type="str">https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1543.003/T1543.003.md</item><item type="str">https://web.archive.org/web/20180331144337/https://www.fireeye.com/blog/threat-research/2018/03/sanny-malware-delivery-method-updated-in-recently-observed-attacks.html</item></references><author type="str">Victor Sergeev, oscd.community, Nasreddine Bencherchali (Nextron Systems)</author><date type="str">2019/10/21</date><modified type="str">2022/11/18</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.privilege_escalation</item><item type="str">attack.t1543.003</item></tags><logsource type="dict"><category type="str">process_creation</category><product type="str">windows</product></logsource><detection type="dict"><selection type="dict"><key name="Image|endswith" type="str">\sc.exe</key><key name="CommandLine|contains|all" type="list"><item type="str">config</item><item type="str">binPath</item></key><key name="CommandLine|contains" type="list"><item type="str">powershell</item><item type="str">cmd </item><item type="str">mshta</item><item type="str">wscript</item><item type="str">cscript</item><item type="str">rundll32</item><item type="str">svchost</item><item type="str">dllhost</item><item type="str">cmd.exe /c</item><item type="str">cmd.exe /k</item><item type="str">cmd.exe /r</item><item type="str">cmd /c</item><item type="str">cmd /k</item><item type="str">cmd /r</item><item type="str">C:\Users\Public</item><item type="str">\Downloads\</item><item type="str">\Desktop\</item><item type="str">\Microsoft\Windows\Start Menu\Programs\Startup\</item><item type="str">C:\Windows\TEMP\</item><item type="str">\AppData\Local\Temp</item></key></selection><condition type="str">selection</condition></detection><fields type="list"><item type="str">CommandLine</item><item type="str">ParentCommandLine</item></fields><falsepositives type="list"><item type="str">Unlikely</item></falsepositives><level type="str">high</level></root>