# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_win_schtasks_env_folder.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Suspicious Schtasks From Env Var Folder</title><id type="str">81325ce1-be01-4250-944f-b4789644556f</id><related type="list"><item type="dict"><id type="str">43f487f0-755f-4c2a-bce7-d6d2eec2fcf8</id><type type="str">derived</type></item></related><status type="str">experimental</status><description type="str">Detects Schtask creations that point to a suspicious folder or an environment variable often used by malware</description><references type="list"><item type="str">https://www.welivesecurity.com/2022/01/18/donot-go-do-not-respawn/</item><item type="str">https://www.joesandbox.com/analysis/514608/0/html#324415FF7D8324231381BAD48A052F85DF04</item></references><author type="str">Florian Roth (Nextron Systems)</author><date type="str">2022/02/21</date><modified type="str">2023/03/14</modified><tags type="list"><item type="str">attack.execution</item><item type="str">attack.t1053.005</item></tags><logsource type="dict"><product type="str">windows</product><category type="str">process_creation</category></logsource><detection type="dict"><selection1_create type="dict"><key name="Image|endswith" type="str">\schtasks.exe</key><key name="CommandLine|contains" type="str"> /create </key></selection1_create><selection1_all_folders type="dict"><key name="CommandLine|contains" type="list"><item type="str">%AppData%</item><item type="str">\AppData\Local\</item><item type="str">\AppData\Roaming\</item><item type="str">%Public%</item><item type="str">\Users\Public</item><item type="str">C:\Windows\Temp</item><item type="str">C:\Perflogs</item></key></selection1_all_folders><selection2_parent type="dict"><key name="ParentCommandLine|endswith" type="str">\svchost.exe -k netsvcs -p -s Schedule</key></selection2_parent><selection2_some_folders type="dict"><key name="CommandLine|contains" type="list"><item type="str">%Public%</item><item type="str">\Users\Public</item><item type="str">C:\Windows\Temp</item><item type="str">C:\Perflogs</item></key></selection2_some_folders><filter_mixed type="list"><item type="dict"><key name="CommandLine|contains" type="list"><item type="str">update_task.xml</item><item type="str">/Create /TN TVInstallRestore /TR</item></key></item><item type="dict"><key name="ParentCommandLine|contains" type="str">unattended.ini</key></item></filter_mixed><filter_avira_install type="dict"><key name="CommandLine|contains|all" type="list"><item type="str">/Create /Xml &quot;C:\Users\</item><item type="str">\AppData\Local\Temp\.CR.</item><item type="str">Avira_Security_Installation.xml</item></key></filter_avira_install><filter_avira_other type="dict"><key name="CommandLine|contains|all" type="list"><item type="str">/Create /F /TN</item><item type="str">/Xml </item><item type="str">\AppData\Local\Temp\is-</item><item type="str">Avira_</item></key><key name="CommandLine|contains" type="list"><item type="str">.tmp\UpdateFallbackTask.xml</item><item type="str">.tmp\WatchdogServiceControlManagerTimeout.xml</item><item type="str">.tmp\SystrayAutostart.xml</item><item type="str">.tmp\MaintenanceTask.xml</item></key></filter_avira_other><filter_klite_codec type="dict"><key name="CommandLine|contains|all" type="list"><item type="str">\AppData\Local\Temp\</item><item type="str">/Create /TN &quot;klcp_update&quot; /XML </item><item type="str">\klcp_update_task.xml</item></key></filter_klite_codec><condition type="str">( all of selection1* or all of selection2* ) and not 1 of filter*</condition></detection><falsepositives type="list"><item type="str">Benign scheduled tasks creations or executions that happen often during software installations</item><item type="str">Software that uses the AppData folder and scheduled tasks to update the software in the AppData folders</item></falsepositives><level type="str">high</level></root>