# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_win_schtasks_reg_loader.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Scheduled Task Executing Payload from Registry</title><id type="str">86588b36-c6d3-465f-9cee-8f9093e07798</id><related type="list"><item type="dict"><id type="str">c4eeeeae-89f4-43a7-8b48-8d1bdfa66c78</id><type type="str">derived</type></item></related><status type="str">experimental</status><description type="str">Detects the creation of a schtasks that potentially executes a payload stored in the Windows Registry using PowerShell.</description><references type="list"><item type="str">https://thedfirreport.com/2022/02/21/qbot-and-zerologon-lead-to-full-domain-compromise/</item></references><author type="str">X__Junior (Nextron Systems), Nasreddine Bencherchali (Nextron Systems)</author><date type="str">2023/07/18</date><tags type="list"><item type="str">attack.execution</item><item type="str">attack.persistence</item><item type="str">attack.t1053.005</item><item type="str">attack.t1059.001</item></tags><logsource type="dict"><product type="str">windows</product><category type="str">process_creation</category></logsource><detection type="dict"><selection_img type="list"><item type="dict"><key name="Image|endswith" type="str">\schtasks.exe</key></item><item type="dict"><OriginalFileName type="str">schtasks.exe</OriginalFileName></item></selection_img><selection_cli_create type="dict"><key name="CommandLine|contains" type="str">/Create</key></selection_cli_create><selection_cli_get type="dict"><key name="CommandLine|contains" type="list"><item type="str">Get-ItemProperty</item><item type="str"> gp </item></key></selection_cli_get><selection_cli_hive type="dict"><key name="CommandLine|contains" type="list"><item type="str">HKCU:</item><item type="str">HKLM:</item><item type="str">registry::</item><item type="str">HKEY_</item></key></selection_cli_hive><filter_main_encoding type="dict"><key name="CommandLine|contains" type="list"><item type="str">FromBase64String</item><item type="str">encodedcommand</item></key></filter_main_encoding><condition type="str">all of selection_* and not 1 of filter_*</condition></detection><falsepositives type="list"><item type="str">Unknown</item></falsepositives><level type="str">medium</level></root>