# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_win_schtasks_susp_pattern.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Suspicious Command Patterns In Scheduled Task Creation</title><id type="str">f2c64357-b1d2-41b7-849f-34d2682c0fad</id><status type="str">experimental</status><description type="str">Detects scheduled task creation using &quot;schtasks&quot; that contain potentially suspicious or uncommon commands</description><references type="list"><item type="str">https://app.any.run/tasks/512c1352-6380-4436-b27d-bb62f0c020d6/</item><item type="str">https://twitter.com/RedDrip7/status/1506480588827467785</item><item type="str">https://www.ncsc.gov.uk/static-assets/documents/malware-analysis-reports/devil-bait/NCSC-MAR-Devil-Bait.pdf</item></references><author type="str">Florian Roth (Nextron Systems)</author><date type="str">2022/02/23</date><modified type="str">2023/05/15</modified><tags type="list"><item type="str">attack.execution</item><item type="str">attack.t1053.005</item></tags><logsource type="dict"><product type="str">windows</product><category type="str">process_creation</category></logsource><detection type="dict"><selection_schtasks type="dict"><key name="Image|endswith" type="str">\schtasks.exe</key><key name="CommandLine|contains" type="str">/Create </key></selection_schtasks><selection_pattern_1 type="dict"><key name="CommandLine|contains" type="list"><item type="str">/sc minute </item><item type="str">/ru system </item></key></selection_pattern_1><selection_pattern_2 type="dict"><key name="CommandLine|contains" type="list"><item type="str">cmd /c</item><item type="str">cmd /k</item><item type="str">cmd /r</item><item type="str">cmd.exe /c </item><item type="str">cmd.exe /k </item><item type="str">cmd.exe /r </item></key></selection_pattern_2><selection_uncommon type="dict"><key name="CommandLine|contains" type="list"><item type="str"> -decode </item><item type="str"> -enc </item><item type="str"> -w hidden </item><item type="str"> bypass </item><item type="str"> IEX</item><item type="str">.DownloadData</item><item type="str">.DownloadFile</item><item type="str">.DownloadString</item><item type="str">/c start /min </item><item type="str">FromBase64String</item><item type="str">mshta http</item><item type="str">mshta.exe http</item></key></selection_uncommon><selection_anomaly_1 type="dict"><key name="CommandLine|contains|all" type="list"><item type="str">\AppData\</item><item type="str">%AppData%</item><item type="str">%Temp%</item><item type="str">%tmp%</item><item type="str">C:\Windows\Temp\</item></key></selection_anomaly_1><selection_anomaly_2 type="dict"><key name="CommandLine|contains|all" type="list"><item type="str">/xml C:\Users\</item><item type="str">cscript</item><item type="str">curl</item><item type="str">wscript</item></key></selection_anomaly_2><condition type="str">selection_schtasks and ( all of selection_pattern_* or selection_uncommon or all of selection_anomaly_* )</condition></detection><falsepositives type="list"><item type="str">Software installers that run from temporary folders and also install scheduled tasks are expected to generate some false positives</item></falsepositives><level type="str">high</level></root>