# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_win_hktl_powersploit_empire_default_schtasks.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">HackTool - Default PowerSploit/Empire Scheduled Task Creation</title><id type="str">56c217c3-2de2-479b-990f-5c109ba8458f</id><status type="str">test</status><description type="str">Detects the creation of a schtask via PowerSploit or Empire Default Configuration.</description><references type="list"><item type="str">https://github.com/0xdeadbeefJERKY/PowerSploit/blob/8690399ef70d2cad10213575ac67e8fa90ddf7c3/Persistence/Persistence.psm1</item><item type="str">https://github.com/EmpireProject/Empire/blob/08cbd274bef78243d7a8ed6443b8364acd1fc48b/lib/modules/powershell/persistence/userland/schtasks.py</item></references><author type="str">Markus Neis, @Karneades</author><date type="str">2018/03/06</date><modified type="str">2023/03/03</modified><tags type="list"><item type="str">attack.execution</item><item type="str">attack.persistence</item><item type="str">attack.privilege_escalation</item><item type="str">attack.s0111</item><item type="str">attack.g0022</item><item type="str">attack.g0060</item><item type="str">car.2013-08-001</item><item type="str">attack.t1053.005</item><item type="str">attack.t1059.001</item></tags><logsource type="dict"><product type="str">windows</product><category type="str">process_creation</category></logsource><detection type="dict"><selection type="dict"><key name="ParentImage|endswith" type="list"><item type="str">\powershell.exe</item><item type="str">\pwsh.exe</item></key><key name="Image|endswith" type="str">\schtasks.exe</key><key name="CommandLine|contains|all" type="list"><item type="str">/Create</item><item type="str">powershell.exe -NonI</item><item type="str">/TN Updater /TR</item></key><key name="CommandLine|contains" type="list"><item type="str">/SC ONLOGON</item><item type="str">/SC DAILY /ST</item><item type="str">/SC ONIDLE</item><item type="str">/SC HOURLY</item></key></selection><condition type="str">selection</condition></detection><falsepositives type="list"><item type="str">Unlikely</item></falsepositives><level type="str">high</level></root>