# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_lnx_at_command.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Scheduled Task/Job At</title><id type="str">d2d642d7-b393-43fe-bae4-e81ed5915c4b</id><status type="str">stable</status><description type="str">Detects the use of at/atd which are utilities that are used to schedule tasks.
They are often abused by adversaries to maintain persistence or to perform task scheduling for initial or recurring execution of malicious code
</description><references type="list"><item type="str">https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1053.002/T1053.002.md</item></references><author type="str">Ömer Günal, oscd.community</author><date type="str">2020/10/06</date><modified type="str">2022/07/07</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1053.002</item></tags><logsource type="dict"><product type="str">linux</product><category type="str">process_creation</category></logsource><detection type="dict"><selection type="dict"><key name="Image|endswith" type="list"><item type="str">/at</item><item type="str">/atd</item></key></selection><condition type="str">selection</condition></detection><falsepositives type="list"><item type="str">Legitimate administration activities</item></falsepositives><level type="str">low</level></root>