# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_win_schtasks_change.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Suspicious Modification Of Scheduled Tasks</title><id type="str">1c0e41cd-21bb-4433-9acc-4a2cd6367b9b</id><related type="list"><item type="dict"><id type="str">614cf376-6651-47c4-9dcc-6b9527f749f4</id><type type="str">similar</type></item></related><status type="str">experimental</status><description type="str">Detects when an attacker tries to modify an already existing scheduled tasks to run from a suspicious location
Attackers can create a simple looking task in order to avoid detection on creation as it&apos;s often the most focused on
Instead they modify the task after creation to include their malicious payload
</description><references type="list"><item type="str">Internal Research</item><item type="str">https://docs.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks</item></references><author type="str">Nasreddine Bencherchali (Nextron Systems)</author><date type="str">2022/07/28</date><modified type="str">2022/11/18</modified><tags type="list"><item type="str">attack.execution</item><item type="str">attack.t1053.005</item></tags><logsource type="dict"><product type="str">windows</product><category type="str">process_creation</category></logsource><detection type="dict"><selection_schtasks type="dict"><key name="Image|endswith" type="str">\schtasks.exe</key><key name="CommandLine|contains|all" type="list"><item type="str"> /Change </item><item type="str"> /TN </item></key></selection_schtasks><selection_susp_locations type="dict"><key name="CommandLine|contains" type="list"><item type="str">\AppData\Local\Temp</item><item type="str">\AppData\Roaming\</item><item type="str">\Users\Public\</item><item type="str">\WINDOWS\Temp\</item><item type="str">\Desktop\</item><item type="str">\Downloads\</item><item type="str">\Temporary Internet</item><item type="str">C:\ProgramData\</item><item type="str">C:\Perflogs\</item><item type="str">%ProgramData%</item><item type="str">%appdata%</item><item type="str">%comspec%</item><item type="str">%localappdata%</item></key></selection_susp_locations><selection_susp_images type="dict"><key name="CommandLine|contains" type="list"><item type="str">regsvr32</item><item type="str">rundll32</item><item type="str">cmd /c </item><item type="str">cmd /k </item><item type="str">cmd /r </item><item type="str">cmd.exe /c </item><item type="str">cmd.exe /k </item><item type="str">cmd.exe /r </item><item type="str">powershell</item><item type="str">mshta</item><item type="str">wscript</item><item type="str">cscript</item><item type="str">certutil</item><item type="str">bitsadmin</item><item type="str">bash.exe</item><item type="str">bash </item><item type="str">scrcons</item><item type="str">wmic </item><item type="str">wmic.exe</item><item type="str">forfiles</item><item type="str">scriptrunner</item><item type="str">hh.exe</item><item type="str">hh </item></key></selection_susp_images><condition type="str">all of selection_*</condition></detection><falsepositives type="list"><item type="str">Unknown</item></falsepositives><level type="str">high</level></root>