# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_win_reg_direct_asep_registry_keys_modification.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Direct Autorun Keys Modification</title><id type="str">24357373-078f-44ed-9ac4-6d334a668a11</id><status type="str">test</status><description type="str">Detects direct modification of autostart extensibility point (ASEP) in registry using reg.exe.</description><references type="list"><item type="str">https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1547.001/T1547.001.md</item></references><author type="str">Victor Sergeev, Daniil Yugoslavskiy, oscd.community</author><date type="str">2019/10/25</date><modified type="str">2022/08/04</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1547.001</item></tags><logsource type="dict"><category type="str">process_creation</category><product type="str">windows</product></logsource><detection type="dict"><selection_1 type="dict"><key name="Image|endswith" type="str">\reg.exe</key><key name="CommandLine|contains" type="str">add</key></selection_1><selection_2 type="dict"><key name="CommandLine|contains" type="list"><item type="str">\software\Microsoft\Windows\CurrentVersion\Run</item><item type="str">\software\Microsoft\Windows NT\CurrentVersion\Winlogon\Userinit</item><item type="str">\software\Microsoft\Windows NT\CurrentVersion\Winlogon\Shell</item><item type="str">\software\Microsoft\Windows NT\CurrentVersion\Windows</item><item type="str">\software\Microsoft\Windows\CurrentVersion\Explorer\User Shell Folders</item><item type="str">\system\CurrentControlSet\Control\SafeBoot\AlternateShell</item></key></selection_2><condition type="str">all of selection_*</condition></detection><fields type="list"><item type="str">CommandLine</item><item type="str">ParentCommandLine</item></fields><falsepositives type="list"><item type="str">Legitimate software automatically (mostly, during installation) sets up autorun keys for legitimate reasons.</item><item type="str">Legitimate administrator sets up autorun keys for legitimate reasons.</item><item type="str">Discord</item></falsepositives><level type="str">medium</level></root>