# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/proc_creation_win_schtasks_system.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Schtasks Creation Or Modification With SYSTEM Privileges</title><id type="str">89ca78fd-b37c-4310-b3d3-81a023f83936</id><status type="str">experimental</status><description type="str">Detects the creation or update of a scheduled task to run with &quot;NT AUTHORITY\SYSTEM&quot; privileges</description><references type="list"><item type="str">https://www.elastic.co/security-labs/exploring-the-qbot-attack-pattern</item><item type="str">https://docs.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks</item></references><author type="str">Nasreddine Bencherchali (Nextron Systems)</author><date type="str">2022/07/28</date><modified type="str">2023/03/14</modified><tags type="list"><item type="str">attack.execution</item><item type="str">attack.persistence</item><item type="str">attack.t1053.005</item></tags><logsource type="dict"><product type="str">windows</product><category type="str">process_creation</category></logsource><detection type="dict"><selection_root type="dict"><key name="Image|endswith" type="str">\schtasks.exe</key><key name="CommandLine|contains" type="list"><item type="str"> /change </item><item type="str"> /create </item></key></selection_root><selection_run type="dict"><key name="CommandLine|contains" type="str">/ru </key></selection_run><selection_user type="dict"><key name="CommandLine|contains" type="list"><item type="str">NT AUT</item><item type="str"> SYSTEM </item></key></selection_user><filter_teamviewer type="dict"><key name="ParentImage|contains|all" type="list"><item type="str">\AppData\Local\Temp\</item><item type="str">TeamViewer_.exe</item></key><key name="Image|endswith" type="str">\schtasks.exe</key><key name="CommandLine|contains" type="str">/TN TVInstallRestore</key></filter_teamviewer><filter_avira type="dict"><key name="CommandLine|contains" type="list"><item type="str">/Create /F /RU System /SC WEEKLY /TN AviraSystemSpeedupVerify /TR </item><item type="str">C:\Program Files (x86)\Avira\System Speedup\setup\avira_speedup_setup.exe</item><item type="str">/VERIFY /VERYSILENT /NOSTART /NODOTNET /NORESTART&quot; /RL HIGHEST</item></key></filter_avira><condition type="str">all of selection_* and not 1 of filter_*</condition></detection><falsepositives type="list"><item type="str">Unknown</item></falsepositives><level type="str">high</level></root>