# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/registry_set_asep_reg_keys_modification_currentversion.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">CurrentVersion Autorun Keys Modification</title><id type="str">20f0ee37-5942-4e45-b7d5-c5b5db9df5cd</id><related type="list"><item type="dict"><id type="str">17f878b8-9968-4578-b814-c4217fc5768c</id><type type="str">obsoletes</type></item></related><status type="str">experimental</status><description type="str">Detects modification of autostart extensibility point (ASEP) in registry.</description><references type="list"><item type="str">https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1547.001/T1547.001.md</item><item type="str">https://docs.microsoft.com/en-us/sysinternals/downloads/autoruns</item><item type="str">https://gist.github.com/GlebSukhodolskiy/0fc5fa5f482903064b448890db1eaf9d</item><item type="str">https://oddvar.moe/2018/03/21/persistence-using-runonceex-hidden-from-autoruns-exe/</item></references><author type="str">Victor Sergeev, Daniil Yugoslavskiy, Gleb Sukhodolskiy, Timur Zinniatullin, oscd.community, Tim Shelton, frack113 (split)</author><date type="str">2019/10/25</date><modified type="str">2022/10/20</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1547.001</item></tags><logsource type="dict"><category type="str">registry_set</category><product type="str">windows</product></logsource><detection type="dict"><current_version_base type="dict"><EventType type="str">SetValue</EventType><key name="TargetObject|contains" type="str">\SOFTWARE\Microsoft\Windows\CurrentVersion</key></current_version_base><current_version_keys type="dict"><key name="TargetObject|contains" type="list"><item type="str">\ShellServiceObjectDelayLoad</item><item type="str">\Run\</item><item type="str">\RunOnce\</item><item type="str">\RunOnceEx\</item><item type="str">\RunServices\</item><item type="str">\RunServicesOnce\</item><item type="str">\Policies\System\Shell</item><item type="str">\Policies\Explorer\Run</item><item type="str">\Group Policy\Scripts\Startup</item><item type="str">\Group Policy\Scripts\Shutdown</item><item type="str">\Group Policy\Scripts\Logon</item><item type="str">\Group Policy\Scripts\Logoff</item><item type="str">\Explorer\ShellServiceObjects</item><item type="str">\Explorer\ShellIconOverlayIdentifiers</item><item type="str">\Explorer\ShellExecuteHooks</item><item type="str">\Explorer\SharedTaskScheduler</item><item type="str">\Explorer\Browser Helper Objects</item><item type="str">\Authentication\PLAP Providers</item><item type="str">\Authentication\Credential Providers</item><item type="str">\Authentication\Credential Provider Filters</item></key></current_version_keys><filter_all type="list"><item type="dict"><Details type="str">(Empty)</Details></item><item type="dict"><key name="TargetObject|endswith" type="str">\NgcFirst\ConsecutiveSwitchCount</key></item><item type="dict"><key name="Image|endswith" type="list"><item type="str">\AppData\Local\Microsoft\OneDrive\Update\OneDriveSetup.exe</item><item type="str">\AppData\Roaming\Spotify\Spotify.exe</item><item type="str">\AppData\Local\WebEx\WebexHost.exe</item></key></item><item type="dict"><Image type="list"><item type="str">C:\WINDOWS\system32\devicecensus.exe</item><item type="str">C:\Windows\system32\winsat.exe</item><item type="str">C:\Program Files\Microsoft OneDrive\StandaloneUpdater\OneDriveSetup.exe</item><item type="str">C:\Program Files\Microsoft OneDrive\Update\OneDriveSetup.exe</item><item type="str">C:\Program Files (x86)\Microsoft OneDrive\Update\OneDriveSetup.exe</item><item type="str">C:\Program Files\KeePass Password Safe 2\ShInstUtil.exe</item><item type="str">C:\Program Files\Everything\Everything.exe</item><item type="str">C:\Program Files (x86)\Microsoft Office\root\integration\integrator.exe</item></Image></item></filter_all><filter_logonui type="dict"><Image type="str">C:\Windows\system32\LogonUI.exe</Image><key name="TargetObject|contains" type="list"><item type="str">\Authentication\Credential Providers\{D6886603-9D2F-4EB2-B667-1971041FA96B}\</item><item type="str">\Authentication\Credential Providers\{BEC09223-B018-416D-A0AC-523971B639F5}\</item><item type="str">\Authentication\Credential Providers\{8AF662BF-65A0-4D0A-A540-A338A999D36F}\</item><item type="str">\Authentication\Credential Providers\{27FBDB57-B613-4AF2-9D7E-4FA7A66C21AD}\</item></key></filter_logonui><filter_edge type="dict"><key name="Image|startswith" type="list"><item type="str">C:\Program Files (x86)\Microsoft\EdgeUpdate\Install\</item><item type="str">C:\Program Files (x86)\Microsoft\EdgeWebView\</item><item type="str">C:\Program Files (x86)\Microsoft\Edge\Application\msedge.exe</item></key></filter_edge><filter_dropbox type="dict"><Image type="str">C:\Windows\system32\regsvr32.exe</Image><key name="TargetObject|contains" type="str">DropboxExt</key><key name="Details|endswith" type="str">A251-47B7-93E1-CDD82E34AF8B}</key></filter_dropbox><filter_opera type="dict"><key name="TargetObject|endswith" type="str">\SOFTWARE\Microsoft\Windows\CurrentVersion\Run\Opera Browser Assistant</key><Details type="str">C:\Program Files\Opera\assistant\browser_assistant.exe</Details></filter_opera><filter_itunes type="dict"><key name="TargetObject|endswith" type="str">\SOFTWARE\Microsoft\Windows\CurrentVersion\Run\iTunesHelper</key><Details type="str">&quot;C:\Program Files\iTunes\iTunesHelper.exe&quot;</Details></filter_itunes><filter_zoom type="dict"><key name="TargetObject|endswith" type="str">\SOFTWARE\Microsoft\Windows\CurrentVersion\RunOnce\zoommsirepair</key><Details type="str">&quot;C:\Program Files\Zoom\bin\installer.exe&quot; /repair</Details></filter_zoom><filter_greenshot type="dict"><key name="TargetObject|endswith" type="str">\SOFTWARE\Microsoft\Windows\CurrentVersion\Run\Greenshot</key><Details type="str">C:\Program Files\Greenshot\Greenshot.exe</Details></filter_greenshot><filter_googledrive1 type="dict"><key name="TargetObject|endswith" type="str">\Software\Microsoft\Windows\CurrentVersion\Run\GoogleDriveFS</key><key name="Details|startswith" type="str">C:\Program Files\Google\Drive File Stream\</key><key name="Details|contains" type="str">\GoogleDriveFS.exe</key></filter_googledrive1><filter_googledrive2 type="dict"><key name="TargetObject|contains" type="str">GoogleDrive</key><Details type="list"><item type="str">{CFE8B367-77A7-41D7-9C90-75D16D7DC6B6}</item><item type="str">{A8E52322-8734-481D-A7E2-27B309EF8D56}</item><item type="str">{C973DA94-CBDF-4E77-81D1-E5B794FBD146}</item><item type="str">{51EF1569-67EE-4AD6-9646-E726C3FFC8A2}</item></Details></filter_googledrive2><filter_onedrive type="dict"><key name="Details|startswith" type="list"><item type="str">C:\Windows\system32\cmd.exe /q /c rmdir /s /q &quot;C:\Users\</item><item type="str">C:\Windows\system32\cmd.exe /q /c del /q &quot;C:\Users\</item></key><key name="Details|contains" type="str">\AppData\Local\Microsoft\OneDrive\</key></filter_onedrive><filter_python type="dict"><key name="TargetObject|contains" type="str">\Microsoft\Windows\CurrentVersion\RunOnce\{</key><key name="Details|contains|all" type="list"><item type="str">\AppData\Local\Package Cache\{</item><item type="str">}\python-</item></key><key name="Details|endswith" type="str">.exe&quot; /burn.runonce</key></filter_python><filter_officeclicktorun type="dict"><key name="Image|startswith" type="list"><item type="str">C:\Program Files\Common Files\Microsoft Shared\ClickToRun\</item><item type="str">C:\Program Files\Common Files\Microsoft Shared\ClickToRun\Updates\</item></key><key name="Image|endswith" type="str">\OfficeClickToRun.exe</key></filter_officeclicktorun><filter_defender type="dict"><Image type="str">C:\Program Files\Windows Defender\MsMpEng.exe</Image></filter_defender><filter_teams type="dict"><key name="Image|endswith" type="str">\Microsoft\Teams\current\Teams.exe</key><key name="Details|contains" type="str">\Microsoft\Teams\Update.exe --processStart </key></filter_teams><filter_ctfmon type="dict"><Image type="str">C:\Windows\system32\userinit.exe</Image><Details type="str">ctfmon.exe /n</Details></filter_ctfmon><filter_AVG type="dict"><key name="Image|startswith" type="str">C:\Program Files\AVG\Antivirus\Setup\</key><Details type="list"><item type="str">&quot;C:\Program Files\AVG\Antivirus\AvLaunch.exe&quot; /gui</item><item type="str">&quot;C:\Program Files (x86)\AVG\Antivirus\AvLaunch.exe&quot; /gui</item><item type="str">{472083B0-C522-11CF-8763-00608CC02F24}</item></Details></filter_AVG><filter_aurora_dashboard type="dict"><key name="Image|endswith" type="list"><item type="str">\aurora-agent-64.exe</item><item type="str">\aurora-agent.exe</item></key><key name="TargetObject|endswith" type="str">\Microsoft\Windows\CurrentVersion\Run\aurora-dashboard</key><Details type="str">C:\Program Files\Aurora-Agent\tools\aurora-dashboard.exe</Details></filter_aurora_dashboard><filter_everything type="dict"><key name="TargetObject|endswith" type="str">\Microsoft\Windows\CurrentVersion\Run\Everything</key><key name="Details|endswith" type="str">\Everything\Everything.exe&quot; -startup</key></filter_everything><condition type="str">all of current_version_* and not 1 of filter_*</condition></detection><fields type="list"><item type="str">SecurityID</item><item type="str">ObjectName</item><item type="str">OldValueType</item><item type="str">NewValueType</item></fields><falsepositives type="list"><item type="str">Legitimate software automatically (mostly, during installation) sets up autorun keys for legitimate reason</item><item type="str">Legitimate administrator sets up autorun keys for legitimate reason</item></falsepositives><level type="str">medium</level></root>