# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/registry_event_susp_download_run_key.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Suspicious Run Key from Download</title><id type="str">9c5037d1-c568-49b3-88c7-9846a5bdc2be</id><status type="str">test</status><description type="str">Detects the suspicious RUN keys created by software located in Download or temporary Outlook/Internet Explorer directories</description><references type="list"><item type="str">https://app.any.run/tasks/c5bef5b7-f484-4c43-9cf3-d5c5c7839def/</item></references><author type="str">Florian Roth (Nextron Systems)</author><date type="str">2019/10/01</date><modified type="str">2021/11/27</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1547.001</item></tags><logsource type="dict"><category type="str">registry_event</category><product type="str">windows</product></logsource><detection type="dict"><selection type="dict"><key name="Image|contains" type="list"><item type="str">\Downloads\</item><item type="str">\Temporary Internet Files\Content.Outlook\</item><item type="str">\Local Settings\Temporary Internet Files\</item></key><key name="TargetObject|contains" type="str">\SOFTWARE\Microsoft\Windows\CurrentVersion\Run\</key></selection><condition type="str">selection</condition></detection><falsepositives type="list"><item type="str">Software installers downloaded and used by users</item></falsepositives><level type="str">high</level></root>