# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/registry_set_taskcache_entry.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Scheduled TaskCache Change by Uncommon Program</title><id type="str">4720b7df-40c3-48fd-bbdf-fd4b3c464f0d</id><status type="str">experimental</status><description type="str">Monitor the creation of a new key under &apos;TaskCache&apos; when a new scheduled task is registered by a process that is not svchost.exe, which is suspicious</description><references type="list"><item type="str">https://thedfirreport.com/2021/03/29/sodinokibi-aka-revil-ransomware/</item><item type="str">https://labs.f-secure.com/blog/scheduled-task-tampering/</item></references><author type="str">Syed Hasan (@syedhasan009)</author><date type="str">2021/06/18</date><modified type="str">2022/10/21</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1053</item><item type="str">attack.t1053.005</item></tags><logsource type="dict"><category type="str">registry_set</category><product type="str">windows</product></logsource><detection type="dict"><selection type="dict"><EventType type="str">SetValue</EventType><key name="TargetObject|contains" type="str">SOFTWARE\Microsoft\Windows NT\CurrentVersion\Schedule\TaskCache\</key></selection><filter type="dict"><key name="TargetObject|contains" type="list"><item type="str">Microsoft\Windows\UpdateOrchestrator</item><item type="str">Microsoft\Windows\SoftwareProtectionPlatform\SvcRestartTask\Index</item><item type="str">Microsoft\Windows\Flighting\OneSettings\RefreshCache\Index</item></key></filter><filter_tiworker type="dict"><key name="Image|startswith" type="str">C:\Windows\</key><key name="Image|endswith" type="str">\TiWorker.exe</key></filter_tiworker><filter_svchost type="dict"><Image type="str">C:\WINDOWS\system32\svchost.exe</Image></filter_svchost><filter_ngen type="dict"><key name="Image|startswith" type="str">C:\Windows\Microsoft.NET\Framework</key><key name="Image|endswith" type="str">\ngen.exe</key><key name="TargetObject|contains" type="list"><item type="str">\Microsoft\Windows NT\CurrentVersion\Schedule\TaskCache\Tasks\{B66B135D-DA06-4FC4-95F8-7458E1D10129}</item><item type="str">\Microsoft\Windows NT\CurrentVersion\Schedule\TaskCache\Tree\Microsoft\Windows\.NET Framework\.NET Framework NGEN</item></key></filter_ngen><filter_office_click_to_run type="dict"><Image type="list"><item type="str">C:\Program Files\Microsoft Office\root\Integration\Integrator.exe</item><item type="str">C:\Program Files (x86)\Microsoft Office\root\Integration\Integrator.exe</item></Image></filter_office_click_to_run><filter_msiexec type="dict"><Image type="str">C:\Windows\System32\msiexec.exe</Image></filter_msiexec><filter_dropbox_updater type="dict"><Image type="list"><item type="str">C:\Program Files (x86)\Dropbox\Update\DropboxUpdate.exe</item><item type="str">C:\Program Files\Dropbox\Update\DropboxUpdate.exe</item></Image></filter_dropbox_updater><filter_explorer type="dict"><Image type="str">C:\Windows\explorer.exe</Image><key name="TargetObject|contains" type="str">\Microsoft\Windows NT\CurrentVersion\Schedule\TaskCache\Tree\Microsoft\Windows\PLA\Server Manager Performance Monitor\</key></filter_explorer><filter_system type="dict"><Image type="str">System</Image></filter_system><condition type="str">selection and not 1 of filter*</condition></detection><falsepositives type="list"><item type="str">Unknown</item></falsepositives><level type="str">high</level></root>