# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/registry_set_servicedll_hijack.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">ServiceDll Hijack</title><id type="str">612e47e9-8a59-43a6-b404-f48683f45bd6</id><status type="str">experimental</status><description type="str">Detects changes to the &quot;ServiceDLL&quot; value related to a service in the registry. This is often used as a method of persistence.</description><references type="list"><item type="str">https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1543.003/T1543.003.md#atomic-test-4---tinyturla-backdoor-service-w64time</item><item type="str">https://www.hexacorn.com/blog/2013/09/19/beyond-good-ol-run-key-part-4/</item></references><author type="str">frack113</author><date type="str">2022/02/04</date><modified type="str">2022/09/20</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.privilege_escalation</item><item type="str">attack.t1543.003</item></tags><logsource type="dict"><category type="str">registry_set</category><product type="str">windows</product></logsource><detection type="dict"><selection type="dict"><EventType type="str">SetValue</EventType><key name="TargetObject|startswith" type="str">HKLM\System\CurrentControlSet\Services\</key><key name="TargetObject|endswith" type="str">\Parameters\ServiceDll</key></selection><filter_printextensionmanger type="dict"><Details type="str">C:\Windows\system32\spool\drivers\x64\3\PrintConfig.dll</Details></filter_printextensionmanger><filter_domain_controller type="dict"><Image type="str">C:\Windows\system32\lsass.exe</Image><key name="TargetObject|endswith" type="str">\CurrentControlSet\Services\NTDS\Parameters\ServiceDll</key><Details type="str">%%systemroot%%\system32\ntdsa.dll</Details></filter_domain_controller><filter_poqexec type="dict"><Image type="str">C:\Windows\System32\poqexec.exe</Image></filter_poqexec><condition type="str">selection and not 1 of filter*</condition></detection><falsepositives type="list"><item type="str">Administrative scripts</item><item type="str">Installation of a service</item></falsepositives><level type="str">medium</level></root>