# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/registry_event_susp_lsass_dll_load.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">DLL Load via LSASS</title><id type="str">b3503044-60ce-4bf4-bbcb-e3db98788823</id><status type="str">test</status><description type="str">Detects a method to load DLL via LSASS process using an undocumented Registry key</description><references type="list"><item type="str">https://blog.xpnsec.com/exploring-mimikatz-part-1/</item><item type="str">https://twitter.com/SBousseaden/status/1183745981189427200</item></references><author type="str">Florian Roth (Nextron Systems)</author><date type="str">2019/10/16</date><modified type="str">2022/04/21</modified><tags type="list"><item type="str">attack.execution</item><item type="str">attack.persistence</item><item type="str">attack.t1547.008</item></tags><logsource type="dict"><category type="str">registry_event</category><product type="str">windows</product></logsource><detection type="dict"><selection type="dict"><key name="TargetObject|contains" type="list"><item type="str">\CurrentControlSet\Services\NTDS\DirectoryServiceExtPt</item><item type="str">\CurrentControlSet\Services\NTDS\LsaDbExtPt</item></key></selection><filter_domain_controller type="dict"><Image type="str">C:\Windows\system32\lsass.exe</Image><Details type="list"><item type="str">%%systemroot%%\system32\ntdsa.dll</item><item type="str">%%systemroot%%\system32\lsadb.dll</item></Details></filter_domain_controller><condition type="str">selection and not 1 of filter_*</condition></detection><falsepositives type="list"><item type="str">Unknown</item></falsepositives><level type="str">high</level></root>