# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/registry_set_asep_reg_keys_modification_wow6432node.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Wow6432Node CurrentVersion Autorun Keys Modification</title><id type="str">b29aed60-ebd1-442b-9cb5-16a1d0324adb</id><related type="list"><item type="dict"><id type="str">17f878b8-9968-4578-b814-c4217fc5768c</id><type type="str">obsoletes</type></item></related><status type="str">experimental</status><description type="str">Detects modification of autostart extensibility point (ASEP) in registry.</description><references type="list"><item type="str">https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1547.001/T1547.001.md</item><item type="str">https://docs.microsoft.com/en-us/sysinternals/downloads/autoruns</item><item type="str">https://gist.github.com/GlebSukhodolskiy/0fc5fa5f482903064b448890db1eaf9d</item><item type="str">https://oddvar.moe/2018/03/21/persistence-using-runonceex-hidden-from-autoruns-exe/</item></references><author type="str">Victor Sergeev, Daniil Yugoslavskiy, Gleb Sukhodolskiy, Timur Zinniatullin, oscd.community, Tim Shelton, frack113 (split)</author><date type="str">2019/10/25</date><modified type="str">2023/01/19</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1547.001</item></tags><logsource type="dict"><category type="str">registry_set</category><product type="str">windows</product></logsource><detection type="dict"><selection_wow_current_version_base type="dict"><EventType type="str">SetValue</EventType><key name="TargetObject|contains" type="str">\SOFTWARE\Wow6432Node\Microsoft\Windows\CurrentVersion</key></selection_wow_current_version_base><selection_wow_current_version_keys type="dict"><key name="TargetObject|contains" type="list"><item type="str">\ShellServiceObjectDelayLoad</item><item type="str">\Run\</item><item type="str">\RunOnce\</item><item type="str">\RunOnceEx\</item><item type="str">\RunServices\</item><item type="str">\RunServicesOnce\</item><item type="str">\Explorer\ShellServiceObjects</item><item type="str">\Explorer\ShellIconOverlayIdentifiers</item><item type="str">\Explorer\ShellExecuteHooks</item><item type="str">\Explorer\SharedTaskScheduler</item><item type="str">\Explorer\Browser Helper Objects</item></key></selection_wow_current_version_keys><filter_empty type="dict"><Details type="str">(Empty)</Details></filter_empty><filter_edge type="dict"><key name="Image|contains|all" type="list"><item type="str">C:\Program Files (x86)\Microsoft\EdgeUpdate\Install\{</item><item type="str">\setup.exe</item></key></filter_edge><filter_msoffice1 type="dict"><Image type="str">C:\Program Files\Common Files\Microsoft Shared\ClickToRun\OfficeClickToRun.exe</Image><key name="TargetObject|contains" type="str">\Office\ClickToRun\REGISTRY\MACHINE\Software\Wow6432Node\</key></filter_msoffice1><filter_msoffice2 type="dict"><Image type="list"><item type="str">C:\Program Files\Microsoft Office\root\integration\integrator.exe</item><item type="str">C:\Program Files (x86)\Microsoft Office\root\integration\integrator.exe</item></Image><key name="TargetObject|contains" type="str">\Explorer\Browser Helper Objects\{31D09BA0-12F5-4CCE-BE8A-2923E76605DA}\</key></filter_msoffice2><filter_dropbox type="list"><item type="dict"><key name="Details|endswith" type="str">-A251-47B7-93E1-CDD82E34AF8B}</key></item><item type="dict"><Details type="str">grpconv -o</Details></item><item type="dict"><key name="Details|contains|all" type="list"><item type="str">C:\Program Files</item><item type="str">\Dropbox\Client\Dropbox.exe</item><item type="str"> /systemstartup</item></key></item></filter_dropbox><filter_evernote type="dict"><key name="TargetObject|endswith" type="str">\Explorer\Browser Helper Objects\{92EF2EAD-A7CE-4424-B0DB-499CF856608E}\NoExplorer</key></filter_evernote><filter_dotnet type="dict"><key name="Image|contains" type="str">\windowsdesktop-runtime-</key><key name="TargetObject|endswith" type="list"><item type="str">\WOW6432Node\Microsoft\Windows\CurrentVersion\RunOnce\{e2d1ae32-dd1d-4ad7-a298-10e42e7840fc}</item><item type="str">\WOW6432Node\Microsoft\Windows\CurrentVersion\RunOnce\{7037b699-7382-448c-89a7-4765961d2537}</item></key><key name="Details|startswith" type="str">&quot;C:\ProgramData\Package Cache\</key><key name="Details|endswith" type="str">.exe&quot; /burn.runonce</key></filter_dotnet><filter_office type="dict"><key name="Image|startswith" type="list"><item type="str">C:\Program Files\Common Files\Microsoft Shared\ClickToRun\</item><item type="str">C:\Program Files\Common Files\Microsoft Shared\ClickToRun\Updates\</item></key><key name="Image|endswith" type="str">\OfficeClickToRun.exe</key></filter_office><filter_ms_win_desktop_runtime type="dict"><key name="Details|startswith" type="str">&quot;C:\ProgramData\Package Cache\{d21a4f20-968a-4b0c-bf04-a38da5f06e41}\windowsdesktop-runtime-</key></filter_ms_win_desktop_runtime><filter_vcredist type="dict"><key name="Image|endswith" type="str">\VC_redist.x64.exe</key><key name="Details|endswith" type="str">}\VC_redist.x64.exe&quot; /burn.runonce</key></filter_vcredist><filter_upgrades type="dict"><key name="Image|startswith" type="list"><item type="str">C:\ProgramData\Package Cache</item><item type="str">C:\Windows\Temp\</item></key><key name="Image|contains" type="list"><item type="str">\winsdksetup.exe</item><item type="str">\windowsdesktop-runtime-</item><item type="str">\AspNetCoreSharedFrameworkBundle-</item></key><key name="Details|endswith" type="str"> /burn.runonce</key></filter_upgrades><filter_uninstallers type="dict"><key name="Image|startswith" type="str">C:\Windows\Installer\MSI</key><key name="TargetObject|contains" type="str">\Explorer\Browser Helper Objects</key></filter_uninstallers><filter_msiexec type="dict"><Image type="str">C:\WINDOWS\system32\msiexec.exe</Image><key name="TargetObject|contains" type="str">\SOFTWARE\Wow6432Node\Microsoft\Windows\CurrentVersion\Run\</key></filter_msiexec><condition type="str">all of selection_wow_current_version_* and not 1 of filter_*</condition></detection><fields type="list"><item type="str">SecurityID</item><item type="str">ObjectName</item><item type="str">OldValueType</item><item type="str">NewValueType</item></fields><falsepositives type="list"><item type="str">Legitimate software automatically (mostly, during installation) sets up autorun keys for legitimate reason</item><item type="str">Legitimate administrator sets up autorun keys for legitimate reason</item></falsepositives><level type="str">medium</level></root>