# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/registry_event_apt_leviathan.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Leviathan Registry Key Activity</title><id type="str">70d43542-cd2d-483c-8f30-f16b436fd7db</id><status type="str">test</status><description type="str">Detects registry key used by Leviathan APT in Malaysian focused campaign</description><references type="list"><item type="str">https://www.elastic.co/blog/advanced-techniques-used-in-malaysian-focused-apt-campaign</item></references><author type="str">Aidan Bracher</author><date type="str">2020/07/07</date><modified type="str">2022/10/09</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1547.001</item></tags><logsource type="dict"><category type="str">registry_event</category><product type="str">windows</product></logsource><detection type="dict"><selection type="dict"><TargetObject type="str">HKCU\Software\Microsoft\Windows\CurrentVersion\Run\ntkd</TargetObject></selection><condition type="str">selection</condition></detection><level type="str">critical</level></root>