# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/registry_set_persistence_com_hijacking_susp_locations.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Potential Persistence Via COM Hijacking From Suspicious Locations</title><id type="str">3d968d17-ffa4-4bc0-bfdc-f139de76ce77</id><status type="str">experimental</status><description type="str">Detects potential COM object hijacking where the &quot;Server&quot; (In/Out) is pointing to a supsicious or unsuale location</description><references type="list"><item type="str">https://www.microsoft.com/security/blog/2022/07/27/untangling-knotweed-european-private-sector-offensive-actor-using-0-day-exploits/ (idea)</item></references><author type="str">Nasreddine Bencherchali (Nextron Systems)</author><date type="str">2022/07/28</date><modified type="str">2023/01/10</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1546.015</item></tags><logsource type="dict"><category type="str">registry_set</category><product type="str">windows</product></logsource><detection type="dict"><selection type="dict"><EventType type="str">SetValue</EventType><key name="TargetObject|startswith" type="list"><item type="str">HKCR\CLSID\</item><item type="str">HKEY_CLASSES_ROOT\CLSID\</item><item type="str">HKCU\Software\Classes\CLSID\</item><item type="str">HKEY_CURRENT_USER\Software\Classes\CLSID\</item></key><key name="TargetObject|endswith" type="list"><item type="str">\InprocServer32\(Default)</item><item type="str">\LocalServer32\(Default)</item></key><key name="Details|contains" type="list"><item type="str">\Users\Public\</item><item type="str">\Desktop\</item><item type="str">\AppData\Local\Temp\</item><item type="str">\Downloads\</item><item type="str">\Windows\Temp\</item><item type="str">\Microsoft\Windows\Start Menu\Programs\Startup\</item><item type="str">\System32\spool\drivers\color\</item><item type="str">%temp%</item><item type="str">%tmp%</item><item type="str">%appdata%</item></key></selection><condition type="str">selection</condition></detection><falsepositives type="list"><item type="str">Probable legitimate applications. If you find these please add them to an exclusion list</item></falsepositives><level type="str">high</level></root>