# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/registry_set_vbs_payload_stored.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">VBScript Payload Stored in Registry</title><id type="str">46490193-1b22-4c29-bdd6-5bf63907216f</id><status type="str">experimental</status><description type="str">Detects VBScript content stored into registry keys as seen being used by UNC2452 group</description><references type="list"><item type="str">https://www.microsoft.com/security/blog/2021/03/04/goldmax-goldfinder-sibot-analyzing-nobelium-malware/</item></references><author type="str">Florian Roth (Nextron Systems)</author><date type="str">2021/03/05</date><modified type="str">2022/09/19</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1547.001</item></tags><logsource type="dict"><category type="str">registry_set</category><product type="str">windows</product></logsource><detection type="dict"><selection type="dict"><EventType type="str">SetValue</EventType><key name="TargetObject|contains" type="str">Software\Microsoft\Windows\CurrentVersion</key><key name="Details|contains" type="list"><item type="str">vbscript:</item><item type="str">jscript:</item><item type="str">mshtml,</item><item type="str">RunHTMLApplication</item><item type="str">Execute(</item><item type="str">CreateObject</item><item type="str">window.close</item></key></selection><filter type="dict"><key name="TargetObject|contains" type="str">Software\Microsoft\Windows\CurrentVersion\Run</key></filter><filter_dotnet type="dict"><key name="Image|endswith" type="str">\msiexec.exe</key><key name="TargetObject|contains" type="str">\SOFTWARE\Microsoft\Windows\CurrentVersion\Installer\UserData\</key><key name="Details|contains" type="list"><item type="str">\Microsoft.NET\Primary Interop Assemblies\Microsoft.mshtml.dll</item><item type="str">&lt;\Microsoft.mshtml,fileVersion=</item><item type="str">_mshtml_dll_</item><item type="str">&lt;\Microsoft.mshtml,culture=</item></key></filter_dotnet><condition type="str">selection and not 1 of filter*</condition></detection><falsepositives type="list"><item type="str">Unknown</item></falsepositives><level type="str">high</level></root>