# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/registry_set_susp_run_key_img_folder.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">New RUN Key Pointing to Suspicious Folder</title><id type="str">02ee49e2-e294-4d0f-9278-f5b3212fc588</id><status type="str">experimental</status><description type="str">Detects suspicious new RUN key element pointing to an executable in a suspicious folder</description><references type="list"><item type="str">https://www.fireeye.com/blog/threat-research/2018/08/fin7-pursuing-an-enigmatic-and-evasive-global-criminal-operation.html</item></references><author type="str">Florian Roth (Nextron Systems), Markus Neis, Sander Wiebing</author><date type="str">2018/08/25</date><modified type="str">2022/09/13</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1547.001</item></tags><logsource type="dict"><category type="str">registry_set</category><product type="str">windows</product></logsource><detection type="dict"><selection_target type="dict"><EventType type="str">SetValue</EventType><key name="TargetObject|contains" type="list"><item type="str">\SOFTWARE\Microsoft\Windows\CurrentVersion\Run\</item><item type="str">\SOFTWARE\Microsoft\Windows\CurrentVersion\RunOnce\</item></key></selection_target><selection_details type="list"><item type="dict"><key name="Details|contains" type="list"><item type="str">C:\Windows\Temp\</item><item type="str">C:\$Recycle.bin\</item><item type="str">C:\Temp\</item><item type="str">C:\Users\Public\</item><item type="str">C:\Users\Default\</item><item type="str">C:\Users\Desktop\</item><item type="str">\AppData\Local\Temp\</item><item type="str">%temp%\</item><item type="str">%tmp%\</item></key></item><item type="dict"><key name="Details|startswith" type="list"><item type="str">%Public%\</item><item type="str">wscript</item><item type="str">cscript</item></key></item></selection_details><condition type="str">all of selection_*</condition></detection><fields type="list"><item type="str">Image</item></fields><falsepositives type="list"><item type="str">Software using weird folders for updates</item></falsepositives><level type="str">high</level></root>