# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/registry_set_asep_reg_keys_modification_currentversion_nt.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">CurrentVersion NT Autorun Keys Modification</title><id type="str">cbf93e5d-ca6c-4722-8bea-e9119007c248</id><related type="list"><item type="dict"><id type="str">17f878b8-9968-4578-b814-c4217fc5768c</id><type type="str">obsoletes</type></item></related><status type="str">experimental</status><description type="str">Detects modification of autostart extensibility point (ASEP) in registry.</description><references type="list"><item type="str">https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1547.001/T1547.001.md</item><item type="str">https://docs.microsoft.com/en-us/sysinternals/downloads/autoruns</item><item type="str">https://gist.github.com/GlebSukhodolskiy/0fc5fa5f482903064b448890db1eaf9d</item></references><author type="str">Victor Sergeev, Daniil Yugoslavskiy, Gleb Sukhodolskiy, Timur Zinniatullin, oscd.community, Tim Shelton, frack113 (split)</author><date type="str">2019/10/25</date><modified type="str">2022/07/05</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1547.001</item></tags><logsource type="dict"><category type="str">registry_set</category><product type="str">windows</product></logsource><detection type="dict"><nt_current_version_base type="dict"><EventType type="str">SetValue</EventType><key name="TargetObject|contains" type="str">\SOFTWARE\Microsoft\Windows NT\CurrentVersion</key></nt_current_version_base><nt_current_version type="dict"><key name="TargetObject|contains" type="list"><item type="str">\Winlogon\VmApplet</item><item type="str">\Winlogon\Userinit</item><item type="str">\Winlogon\Taskman</item><item type="str">\Winlogon\Shell</item><item type="str">\Winlogon\GpExtensions</item><item type="str">\Winlogon\AppSetup</item><item type="str">\Winlogon\AlternateShells\AvailableShells</item><item type="str">\Windows\IconServiceLib</item><item type="str">\Windows\Appinit_Dlls</item><item type="str">\Image File Execution Options</item><item type="str">\Font Drivers</item><item type="str">\Drivers32</item><item type="str">\Windows\Run</item><item type="str">\Windows\Load</item></key></nt_current_version><filter_empty type="dict"><Details type="str">(Empty)</Details></filter_empty><filter_legitimate_subkey type="dict"><key name="TargetObject|contains" type="str">\Image File Execution Options\</key><key name="TargetObject|endswith" type="list"><item type="str">\DisableExceptionChainValidation</item><item type="str">\MitigationOptions</item></key></filter_legitimate_subkey><filter_edge type="dict"><key name="Image|startswith" type="str">C:\Program Files (x86)\Microsoft\Temp\</key><key name="Image|endswith" type="str">\MicrosoftEdgeUpdate.exe</key></filter_edge><filter_msoffice type="list"><item type="dict"><key name="TargetObject|contains" type="list"><item type="str">\ClickToRunStore\HKLM\SOFTWARE\Microsoft\Windows NT\CurrentVersion\</item><item type="str">\ClickToRun\REGISTRY\MACHINE\Software\Microsoft\Windows NT\CurrentVersion\</item></key></item><item type="dict"><Image type="list"><item type="str">C:\Program Files\Microsoft Office\root\integration\integrator.exe</item><item type="str">C:\Program Files (x86)\Microsoft Office\root\integration\integrator.exe</item></Image></item></filter_msoffice><filter_officeclicktorun type="dict"><key name="Image|startswith" type="list"><item type="str">C:\Program Files\Common Files\Microsoft Shared\ClickToRun\</item><item type="str">C:\Program Files\Common Files\Microsoft Shared\ClickToRun\Updates\</item></key><key name="Image|endswith" type="str">\OfficeClickToRun.exe</key></filter_officeclicktorun><filter_security_extension_dc type="dict"><Image type="str">C:\Windows\system32\svchost.exe</Image><key name="TargetObject|contains" type="list"><item type="str">\Winlogon\GPExtensions\{827D319E-6EAC-11D2-A4EA-00C04F79F83A}\PreviousPolicyAreas</item><item type="str">\Winlogon\GPExtensions\{827D319E-6EAC-11D2-A4EA-00C04F79F83A}\MaxNoGPOListChangesInterval</item></key><Details type="list"><item type="str">DWORD (0x00000009)</item><item type="str">DWORD (0x000003c0)</item></Details></filter_security_extension_dc><filter_ngen type="dict"><key name="Image|startswith" type="str">C:\Windows\Microsoft.NET\Framework</key><key name="Image|endswith" type="str">\ngen.exe</key></filter_ngen><filter_onedrive type="dict"><key name="Image|endswith" type="str">\AppData\Local\Microsoft\OneDrive\StandaloneUpdater\OneDriveSetup.exe</key><key name="TargetObject|endswith" type="str">\Microsoft\Windows\CurrentVersion\RunOnce\Delete Cached Update Binary</key><key name="Details|startswith" type="str">C:\Windows\system32\cmd.exe /q /c del /q &quot;C:\Users\</key><key name="Details|endswith" type="str">\AppData\Local\Microsoft\OneDrive\Update\OneDriveSetup.exe&quot;</key></filter_onedrive><condition type="str">nt_current_version_base and nt_current_version and not 1 of filter_*</condition></detection><fields type="list"><item type="str">SecurityID</item><item type="str">ObjectName</item><item type="str">OldValueType</item><item type="str">NewValueType</item></fields><falsepositives type="list"><item type="str">Legitimate software automatically (mostly, during installation) sets up autorun keys for legitimate reason</item><item type="str">Legitimate administrator sets up autorun keys for legitimate reason</item></falsepositives><level type="str">medium</level></root>