# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/registry_set_timeproviders_dllname.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Set TimeProviders DllName</title><id type="str">e88a6ddc-74f7-463b-9b26-f69fc0d2ce85</id><status type="str">experimental</status><description type="str">Detects processes setting a new DLL in DllName in under HKEY_LOCAL_MACHINE\SYSTEM\CurrentControlSet\Services\W32Time\TimeProvider.
Adversaries may abuse time providers to execute DLLs when the system boots.
The Windows Time service (W32Time) enables time synchronization across and within domains.
</description><references type="list"><item type="str">https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1547.003/T1547.003.md</item></references><author type="str">frack113</author><date type="str">2022/06/19</date><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.privilege_escalation</item><item type="str">attack.t1547.003</item></tags><logsource type="dict"><category type="str">registry_set</category><product type="str">windows</product></logsource><detection type="dict"><selection type="dict"><EventType type="str">SetValue</EventType><key name="TargetObject|startswith" type="str">HKLM\System\CurrentControlSet\Services\W32Time\TimeProviders</key><key name="TargetObject|endswith" type="str">DllName</key></selection><filter type="dict"><Details type="str">C:\Windows\SYSTEM32\w32time.DLL</Details></filter><condition type="str">selection and not filter</condition></detection><falsepositives type="list"><item type="str">Unknown</item></falsepositives><level type="str">high</level></root>