# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/registry_event_runkey_winekey.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">WINEKEY Registry Modification</title><id type="str">b98968aa-dbc0-4a9c-ac35-108363cbf8d5</id><status type="str">test</status><description type="str">Detects potential malicious modification of run keys by winekey or team9 backdoor</description><references type="list"><item type="str">https://www.fireeye.com/blog/threat-research/2020/10/kegtap-and-singlemalt-with-a-ransomware-chaser.html</item></references><author type="str">omkar72</author><date type="str">2020/10/30</date><modified type="str">2021/11/27</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1547</item></tags><logsource type="dict"><category type="str">registry_event</category><product type="str">windows</product></logsource><detection type="dict"><selection type="dict"><key name="TargetObject|endswith" type="str">Software\Microsoft\Windows\CurrentVersion\Run\Backup Mgr</key></selection><condition type="str">selection</condition></detection><fields type="list"><item type="str">ComputerName</item><item type="str">Image</item><item type="str">EventType</item><item type="str">TargetObject</item></fields><falsepositives type="list"><item type="str">Unknown</item></falsepositives><level type="str">high</level></root>