# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/registry_set_powershell_in_run_keys.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Suspicious Powershell In Registry Run Keys</title><id type="str">8d85cf08-bf97-4260-ba49-986a2a65129c</id><status type="str">experimental</status><description type="str">Detects potential PowerShell commands or code within registry run keys</description><references type="list"><item type="str">https://github.com/frack113/atomic-red-team/blob/a9051c38de8a5320b31c7039efcbd3b56cf2d65a/atomics/T1547.001/T1547.001.md#atomic-test-9---systembc-malware-as-a-service-registry</item><item type="str">https://www.trendmicro.com/en_us/research/22/j/lv-ransomware-exploits-proxyshell-in-attack.html</item></references><author type="str">frack113, Florian Roth (Nextron Systems)</author><date type="str">2022/03/17</date><modified type="str">2023/01/19</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1547.001</item></tags><logsource type="dict"><category type="str">registry_set</category><product type="str">windows</product></logsource><detection type="dict"><selection type="dict"><EventType type="str">SetValue</EventType><key name="TargetObject|contains" type="str">\Software\Microsoft\Windows\CurrentVersion\Run</key><key name="Details|contains" type="list"><item type="str">powershell</item><item type="str">pwsh </item><item type="str">FromBase64String</item><item type="str">.DownloadFile(</item><item type="str">.DownloadString(</item><item type="str"> -w hidden </item><item type="str"> -w 1 </item><item type="str">-windowstyle hidden</item><item type="str">-window hidden</item><item type="str"> -nop </item><item type="str"> -encodedcommand </item><item type="str">-ExecutionPolicy Bypass</item><item type="str">Invoke-Expression</item><item type="str">IEX (</item><item type="str">Invoke-Command</item><item type="str">ICM -</item><item type="str">Invoke-WebRequest</item><item type="str">IWR </item><item type="str"> -noni </item><item type="str"> -noninteractive </item></key></selection><condition type="str">selection</condition></detection><falsepositives type="list"><item type="str">Legitimate admin or third party scripts. Baseline according to your environment</item></falsepositives><level type="str">medium</level></root>