# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/registry_set_asep_reg_keys_modification_winsock2.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">WinSock2 Autorun Keys Modification</title><id type="str">d6c2ce7e-afb5-4337-9ca4-4b5254ed0565</id><related type="list"><item type="dict"><id type="str">17f878b8-9968-4578-b814-c4217fc5768c</id><type type="str">derived</type></item></related><status type="str">test</status><description type="str">Detects modification of autostart extensibility point (ASEP) in registry.</description><references type="list"><item type="str">https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1547.001/T1547.001.md</item><item type="str">https://docs.microsoft.com/en-us/sysinternals/downloads/autoruns</item><item type="str">https://gist.github.com/GlebSukhodolskiy/0fc5fa5f482903064b448890db1eaf9d</item></references><author type="str">Victor Sergeev, Daniil Yugoslavskiy, Gleb Sukhodolskiy, Timur Zinniatullin, oscd.community, Tim Shelton, frack113 (split)</author><date type="str">2019/10/25</date><modified type="str">2022/03/26</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1547.001</item></tags><logsource type="dict"><category type="str">registry_set</category><product type="str">windows</product></logsource><detection type="dict"><winsock_parameters_base type="dict"><EventType type="str">SetValue</EventType><key name="TargetObject|contains" type="str">\System\CurrentControlSet\Services\WinSock2\Parameters</key></winsock_parameters_base><winsock_parameters type="dict"><key name="TargetObject|contains" type="list"><item type="str">\Protocol_Catalog9\Catalog_Entries</item><item type="str">\NameSpace_Catalog5\Catalog_Entries</item></key></winsock_parameters><filter type="list"><item type="dict"><Details type="str">(Empty)</Details></item><item type="dict"><Image type="str">C:\Windows\System32\MsiExec.exe</Image></item><item type="dict"><Image type="str">C:\Windows\syswow64\MsiExec.exe</Image></item></filter><condition type="str">winsock_parameters_base and winsock_parameters and not filter</condition></detection><fields type="list"><item type="str">SecurityID</item><item type="str">ObjectName</item><item type="str">OldValueType</item><item type="str">NewValueType</item></fields><falsepositives type="list"><item type="str">Legitimate software automatically (mostly, during installation) sets up autorun keys for legitimate reason</item><item type="str">Legitimate administrator sets up autorun keys for legitimate reason</item></falsepositives><level type="str">medium</level></root>