# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/rpc_firewall_atsvc_lateral_movement.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Remote Schedule Task Lateral Movement via ATSvc</title><id type="str">0fcd1c79-4eeb-4746-aba9-1b458f7a79cb</id><status type="str">test</status><description type="str">Detects remote RPC calls to create or execute a scheduled task via ATSvc</description><references type="list"><item type="str">https://docs.microsoft.com/en-us/openspecs/windows_protocols/ms-tsch/d1058a28-7e02-4948-8b8d-4a347fa64931</item><item type="str">https://github.com/jsecurity101/MSRPC-to-ATTACK/blob/ddd4608fe8684fcf2fcf9b48c5f0b3c28097f8a3/documents/MS-TSCH.md</item><item type="str">https://github.com/zeronetworks/rpcfirewall</item><item type="str">https://zeronetworks.com/blog/stopping-lateral-movement-via-the-rpc-firewall/</item></references><author type="str">Sagie Dulce, Dekel Paz</author><date type="str">2022/01/01</date><modified type="str">2022/01/01</modified><tags type="list"><item type="str">attack.lateral_movement</item><item type="str">attack.t1053</item><item type="str">attack.t1053.002</item></tags><logsource type="dict"><product type="str">rpc_firewall</product><category type="str">application</category><definition type="str">Requirements: install and apply the RPC Firewall to all processes with &quot;audit:true action:block uuid:1ff70682-0a51-30e8-076d-740be8cee98b&quot;</definition></logsource><detection type="dict"><selection type="dict"><EventLog type="str">RPCFW</EventLog><EventID type="int">3</EventID><InterfaceUuid type="str">1ff70682-0a51-30e8-076d-740be8cee98b</InterfaceUuid><OpNum type="list"><item type="int">0</item><item type="int">1</item></OpNum></selection><condition type="str">selection</condition></detection><falsepositives type="list"><item type="str">Unknown</item></falsepositives><level type="str">high</level></root>