# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/win_security_susp_scheduled_task_update.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Suspicious Scheduled Task Update</title><id type="str">614cf376-6651-47c4-9dcc-6b9527f749f4</id><related type="list"><item type="dict"><id type="str">1c0e41cd-21bb-4433-9acc-4a2cd6367b9b</id><type type="str">similar</type></item></related><status type="str">experimental</status><description type="str">Detects update to a scheduled task event that contain suspicious keywords.</description><references type="list"><item type="str">https://learn.microsoft.com/en-us/windows/security/threat-protection/auditing/event-4698</item></references><author type="str">Nasreddine Bencherchali (Nextron Systems)</author><date type="str">2022/12/05</date><tags type="list"><item type="str">attack.execution</item><item type="str">attack.privilege_escalation</item><item type="str">attack.persistence</item><item type="str">attack.t1053.005</item></tags><logsource type="dict"><product type="str">windows</product><service type="str">security</service><definition type="str">The Advanced Audit Policy setting Object Access &gt; Audit Other Object Access Events has to be configured to allow this detection. We also recommend extracting the Command field from the embedded XML in the event data.</definition></logsource><detection type="dict"><selection_eid type="dict"><EventID type="int">4702</EventID></selection_eid><selection_paths type="dict"><key name="TaskContentNew|contains" type="list"><item type="str">\AppData\Local\Temp\</item><item type="str">\AppData\Roaming\</item><item type="str">\Users\Public\</item><item type="str">\WINDOWS\Temp\</item><item type="str">C:\Temp\</item><item type="str">\Desktop\</item><item type="str">\Downloads\</item><item type="str">\Temporary Internet</item><item type="str">C:\ProgramData\</item><item type="str">C:\Perflogs\</item></key></selection_paths><selection_commands type="dict"><key name="TaskContentNew|contains" type="list"><item type="str">regsvr32</item><item type="str">rundll32</item><item type="str">cmd.exe&lt;/Command&gt;</item><item type="str">cmd&lt;/Command&gt;</item><item type="str">&lt;Arguments&gt;/c </item><item type="str">&lt;Arguments&gt;/k </item><item type="str">&lt;Arguments&gt;/r </item><item type="str">powershell</item><item type="str">pwsh</item><item type="str">mshta</item><item type="str">wscript</item><item type="str">cscript</item><item type="str">certutil</item><item type="str">bitsadmin</item><item type="str">bash.exe</item><item type="str">bash </item><item type="str">scrcons</item><item type="str">wmic </item><item type="str">wmic.exe</item><item type="str">forfiles</item><item type="str">scriptrunner</item><item type="str">hh.exe</item></key></selection_commands><condition type="str">all of selection_*</condition></detection><falsepositives type="list"><item type="str">Unknown</item></falsepositives><level type="str">high</level></root>