# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/win_taskscheduler_execution_from_susp_locations.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Scheduled Task Executed From A Suspicious Location</title><id type="str">424273ea-7cf8-43a6-b712-375f925e481f</id><status type="str">experimental</status><description type="str">Detects the execution of Scheduled Tasks where the Program being run is located in a suspicious location or it&apos;s an unusale program to be run from a Scheduled Task</description><references type="list"><item type="str">Internal Research</item></references><author type="str">Nasreddine Bencherchali (Nextron Systems)</author><date type="str">2022/12/05</date><modified type="str">2023/02/07</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.t1053.005</item></tags><logsource type="dict"><product type="str">windows</product><service type="str">taskscheduler</service><definition type="str">Requirements: The &quot;Microsoft-Windows-TaskScheduler/Operational&quot; is disabled by default and needs to be enabled in order for this detection to trigger</definition></logsource><detection type="dict"><selection type="dict"><EventID type="int">129</EventID><key name="Path|contains" type="list"><item type="str">C:\Windows\Temp\</item><item type="str">\AppData\Local\Temp\</item><item type="str">\Desktop\</item><item type="str">\Downloads\</item><item type="str">\Users\Public\</item><item type="str">C:\Temp\</item></key></selection><condition type="str">selection</condition></detection><falsepositives type="list"><item type="str">Unknown</item></falsepositives><level type="str">medium</level></root>