# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/win_security_susp_scheduled_task_delete_or_disable.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Important Scheduled Task Deleted/Disabled</title><id type="str">7595ba94-cf3b-4471-aa03-4f6baa9e5fad</id><related type="list"><item type="dict"><id type="str">dbc1f800-0fe0-4bc0-9c66-292c2abe3f78</id><type type="str">similar</type></item><item type="dict"><id type="str">9ac94dc8-9042-493c-ba45-3b5e7c86b980</id><type type="str">similar</type></item><item type="dict"><id type="str">9e3cb244-bdb8-4632-8c90-6079c8f4f16d</id><type type="str">similar</type></item></related><status type="str">experimental</status><description type="str">Detects when adversaries stop services or processes by deleting or disabling their respective scheduled tasks in order to conduct data destructive activities</description><references type="list"><item type="str">https://learn.microsoft.com/en-us/windows/security/threat-protection/auditing/event-4699</item><item type="str">https://learn.microsoft.com/en-us/windows/security/threat-protection/auditing/event-4701</item></references><author type="str">Nasreddine Bencherchali (Nextron Systems)</author><date type="str">2022/12/05</date><modified type="str">2023/03/13</modified><tags type="list"><item type="str">attack.execution</item><item type="str">attack.privilege_escalation</item><item type="str">attack.persistence</item><item type="str">attack.t1053.005</item></tags><logsource type="dict"><product type="str">windows</product><service type="str">security</service><definition type="str">The Advanced Audit Policy setting Object Access &gt; Audit Other Object Access Events has to be configured to allow this detection. We also recommend extracting the Command field from the embedded XML in the event data.</definition></logsource><detection type="dict"><selection type="dict"><EventID type="list"><item type="int">4699</item><item type="int">4701</item></EventID><key name="TaskName|contains" type="list"><item type="str">\Windows\SystemRestore\SR</item><item type="str">\Windows\Windows Defender\</item><item type="str">\Windows\BitLocker</item><item type="str">\Windows\WindowsBackup\</item><item type="str">\Windows\WindowsUpdate\</item><item type="str">\Windows\UpdateOrchestrator\Schedule</item><item type="str">\Windows\ExploitGuard</item></key></selection><filter_sys_username type="dict"><EventID type="int">4699</EventID><key name="SubjectUserName|endswith" type="str">$</key><key name="TaskName|contains" type="str">\Windows\Windows Defender\</key></filter_sys_username><condition type="str">selection and not 1 of filter_*</condition></detection><falsepositives type="list"><item type="str">Unknown</item></falsepositives><level type="str">high</level></root>