# ravin_correlation_rules_git/sigma_rules/convertor/sigma_rules_xml/win_security_gpo_scheduledtasks.xml
<?xml version="1.0" encoding="UTF-8" ?><root><title type="str">Persistence and Execution at Scale via GPO Scheduled Task</title><id type="str">a8f29a7b-b137-4446-80a0-b804272f3da2</id><status type="str">test</status><description type="str">Detect lateral movement using GPO scheduled task, usually used to deploy ransomware at scale</description><references type="list"><item type="str">https://twitter.com/menasec1/status/1106899890377052160</item><item type="str">https://www.secureworks.com/blog/ransomware-as-a-distraction</item></references><author type="str">Samir Bousseaden</author><date type="str">2019/04/03</date><modified type="str">2021/11/27</modified><tags type="list"><item type="str">attack.persistence</item><item type="str">attack.lateral_movement</item><item type="str">attack.t1053.005</item></tags><logsource type="dict"><product type="str">windows</product><service type="str">security</service><definition type="str">The advanced audit policy setting &quot;Object Access &gt; Audit Detailed File Share&quot; must be configured for Success/Failure</definition></logsource><detection type="dict"><selection type="dict"><EventID type="int">5145</EventID><ShareName type="str">\\\\\*\\SYSVOL</ShareName><key name="RelativeTargetName|endswith" type="str">ScheduledTasks.xml</key><key name="Accesses|contains" type="list"><item type="str">WriteData</item><item type="str">%%4417</item></key></selection><condition type="str">selection</condition></detection><falsepositives type="list"><item type="str">If the source IP is not localhost then it&apos;s super suspicious, better to monitor both local and remote changes to GPO scheduledtasks</item></falsepositives><level type="str">high</level></root>